# squashfs-root/lib/udev/rules.d/50-apport.rules
# firmware dumps from iwlmvm iwlwifi driver
DRIVER=="iwlwifi", ACTION=="change", ENV{EVENT}=="error_dump", RUN+="/usr/share/apport/iwlwifi_error_dump $env{DEVPATH}"
