# squashfs-root/etc/init.d/whoopsie
#!/lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          whoopsie
# Required-Start:    $network $syslog $remote_fs
# Required-Stop:     $network $syslog $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      1
# Short-Description: crash report submission daemon
# Description:       Debian init script to start the crash report
#                    submission daemon
### END INIT INFO
DAEMON=/bin/sh -c 'export CRASH_DB_URL=https://daisy.ubuntu.com; exec whoopsie -f'

