# squashfs-root/snap/core/7270/etc/systemd/system/swapfile.service
[Unit]
Description=Create and use a swap file

[Service]
ExecStart=/usr/bin/mkswapfile

[Install]
WantedBy=sysinit.target
