# squashfs-root/snap/core/7270/lib/systemd/system/fwupdate-cleanup.service
[Unit]
Description=fwupdate housekeeping

[Service]
PrivateTmp=true
Type=oneshot
ExecStart=/usr/lib/fwupdate/cleanup
