# squashfs-root/snap/core/7270/lib/systemd/system/systemd-networkd-wait-online.path
[Unit]
Description=Trigger systemd-networkd-wait-online if netplan updates

[Path]
PathChanged=/run/systemd/generator/netplan.stamp
