# squashfs-root/snap/core/7270/lib/systemd/system/ubuntu-fan.service
[Unit]
Description=Ubuntu FAN network setup
After=network-online.target
Before=lxc.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/ubuntu-fan/fan-net start
ExecStop=/usr/lib/ubuntu-fan/fan-net stop

[Install]
WantedBy=multi-user.target
