# squashfs-root/usr/share/upstart/sessions/session-migration.conf
description "Session Migration"
author "Iain Lane <iain.lane@canonical.com>"

start on (started dbus and starting xsession-init)

task

exec session-migration
