# squashfs-root/snap/core/7270/lib/systemd/system/systemd-resolved.service.d/resolvconf.conf
# When resolved is in use, it must be brought up before we consider networking
# available because otherwise there is a window where DNS resolution doesn't
# work.
[Unit]
Before=network-online.target
