# squashfs-root/snap/core18/1066/lib/systemd/system/core18.start-snapd.service
[Unit]
Description=Start the snapd services from the snapd snap
RequiresMountsFor=/run

[Service]
ExecStart=/usr/lib/core18/run-snapd-from-snap start
Type=oneshot
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
