# squashfs-root/snap/core18/1066/lib/systemd/system/snapd.seeded.service
[Unit]
Description=Wait until snapd is fully seeded (core18)
After=core18.start-snapd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/snap wait system seed.loaded
RemainAfterExit=true

[Install]
WantedBy=multi-user.target cloud-final.service
