# squashfs-root/usr/lib/systemd/user/at-spi-dbus-bus.service
[Unit]
Description=Accessibility services bus

[Service]
Type=dbus
BusName=org.a11y.Bus
ExecStart=/usr/lib/at-spi2-core/at-spi-bus-launcher
