# squashfs-root/usr/lib/systemd/user/colord-session.service
[Unit]
Description=Color management helper

[Service]
Type=dbus
BusName=org.freedesktop.ColorHelper
ExecStart=/usr/lib/colord/colord-session
