# squashfs-root/usr/lib/systemd/user/flatpak-portal.service
[Unit]
Description=flatpak portal

[Service]
BusName=org.freedesktop.portal.Flatpak
ExecStart=/usr/lib/flatpak/flatpak-portal
Type=dbus
