# squashfs-root/usr/lib/systemd/user/gnome-keyring.service
[Unit]
Description=Start gnome-keyring for the Secrets Service, and PKCS #11
Before=graphical-session-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/gnome-keyring-daemon --start --components pkcs11,secrets

[Install]
WantedBy=graphical-session-pre.target
