# squashfs-root/usr/lib/systemd/user/gnome-terminal-server.service
[Unit]
Description=GNOME Terminal Server
[Service]
KillMode=process
Type=dbus
BusName=org.gnome.Terminal
ExecStart=/usr/lib/gnome-terminal/gnome-terminal-server
