# squashfs-root/usr/lib/systemd/user/gvfs-metadata.service
[Unit]
Description=Virtual filesystem metadata service

[Service]
ExecStart=/usr/lib/gvfs/gvfsd-metadata
Type=dbus
BusName=org.gtk.vfs.Metadata
