# squashfs-root/usr/lib/systemd/user/obex.service
[Unit]
Description=Bluetooth OBEX service

[Service]
Type=dbus
BusName=org.bluez.obex
ExecStart=/usr/lib/bluetooth/obexd

[Install]
Alias=dbus-org.bluez.obex.service
