# squashfs-root/usr/lib/systemd/user/unicast-local-avahi.path
[Unit]
Description=Path trigger for Avahi .local domain notifications
PartOf=graphical-session.target

[Path]
PathChanged=/run/avahi-daemon/disabled-for-unicast-local
