# squashfs-root/usr/lib/systemd/user/unicast-local-avahi.service
[Unit]
Description=Notification regarding avahi disabled due to .local domain

[Service]
ExecStart=/usr/lib/update-notifier/local-avahi-notification
Type=oneshot
