# squashfs-root/usr/lib/systemd/user/update-notifier-livepatch.path
[Unit]
Description=Path trigger for Livepatch status notifications
PartOf=graphical-session.target

[Path]
PathChanged=/var/snap/canonical-livepatch/current/status
