# squashfs-root/usr/lib/systemd/user/update-notifier-livepatch.service
[Unit]
Description=Livepatch status notification

[Service]
ExecStart=/usr/lib/update-notifier/livepatch-notification
Type=oneshot
