# squashfs-root/usr/lib/systemd/user/update-notifier-release.service
[Unit]
Description=Notification regarding a new release of Ubuntu
PartOf=graphical-session.target

[Service]
ExecStart=/usr/lib/ubuntu-release-upgrader/check-new-release-gtk
Type=oneshot
