# squashfs-root/usr/lib/systemd/user/update-notifier-release.path
[Unit]
Description=Path trigger for new release of Ubuntu notifications
PartOf=graphical-session.target

[Path]
PathChanged=/var/lib/ubuntu-release-upgrader/release-upgrade-available
