# squashfs-root/usr/lib/systemd/user/update-notifier-crash.path
[Unit]
Description=Path trigger for Apport crash notifications
PartOf=graphical-session.target

[Path]
PathChanged=/var/crash/
