# squashfs-root/etc/network/if-up.d/update
#!/bin/sh
apt-get update
#!/bin/sh
apt-get update
