# squashfs-root/snap/core18/1066/lib/systemd/system/systemd-tmpfiles-setup.service.d/fixup-timesyncd.conf
[Service]
ExecStartPost=/bin/sh -c "[ ! -L /var/lib/systemd/timesync ] && rm -rf /var/lib/systemd/timesync || true"
