# squashfs-root/snap/core/7270/lib/udev/rules.d/40-bridge-network-interface.rules
ACTION=="add", SUBSYSTEM=="net", RUN+="bridge-network-interface"
