# squashfs-root/snap/core/7270/lib/udev/rules.d/85-keyboard-configuration.rules
# Set up a single virtual console.
SUBSYSTEM=="tty", KERNEL=="tty[1-6]", RUN+="console-setup-tty /dev/%k"

# Framebuffer console created.  Set up the font on all virtual consoles
# again.  This races with Plymouth if fbcon is built into the kernel; but
# the worst case should simply be that fonts don't get set up properly.
SUBSYSTEM=="graphics", KERNEL=="fbcon", RUN+="console-setup-tty fbcon"
