# squashfs-root/usr/share/gvfs/mounts/afc.mount
[Mount]
Type=afc
Exec=/usr/lib/gvfs/gvfsd-afc
AutoMount=false
Scheme=afc
DefaultPort=1

