# squashfs-root/usr/share/gvfs/mounts/afp-browse.mount
[Mount]
Type=afp-server
Exec=/usr/lib/gvfs/gvfsd-afp-browse
AutoMount=false
Scheme=afp
DefaultPort=548
HostnameIsInetAddress=true

