# squashfs-root/usr/share/gvfs/mounts/afp.mount
[Mount]
Type=afp-volume
Exec=/usr/lib/gvfs/gvfsd-afp
AutoMount=false
Scheme=afp
DefaultPort=548
HostnameIsInetAddress=true