# squashfs-root/usr/share/gvfs/mounts/archive.mount
[Mount]
Type=archive
Exec=/usr/lib/gvfs/gvfsd-archive
AutoMount=false
Scheme=archive
