# squashfs-root/usr/share/gvfs/mounts/burn.mount
[Mount]
Type=burn
Exec=/usr/lib/gvfs/gvfsd-burn
AutoMount=true

