# squashfs-root/usr/share/gvfs/mounts/cdda.mount
[Mount]
Type=cdda
Exec=/usr/lib/gvfs/gvfsd-cdda
AutoMount=false
