# squashfs-root/usr/share/gvfs/mounts/computer.mount
[Mount]
Type=computer
Exec=/usr/lib/gvfs/gvfsd-computer
AutoMount=true

