# squashfs-root/usr/share/gvfs/mounts/dav+sd.mount
[Mount]
Type=dav+sd;davs+sd
Exec=/usr/lib/gvfs/gvfsd-dav
AutoMount=false
