# squashfs-root/usr/share/gvfs/mounts/dav.mount
[Mount]
Type=dav;davs
Exec=/usr/lib/gvfs/gvfsd-dav
AutoMount=false
