# squashfs-root/usr/share/gvfs/mounts/dns-sd.mount
[Mount]
Type=dns-sd
Exec=/usr/lib/gvfs/gvfsd-dnssd
AutoMount=true
DBusName=org.gtk.vfs.mountpoint_dnssd
