# squashfs-root/usr/share/gvfs/mounts/ftp.mount
[Mount]
Type=ftp
Exec=/usr/lib/gvfs/gvfsd-ftp
AutoMount=false
Scheme=ftp
DefaultPort=21
HostnameIsInetAddress=true
