# squashfs-root/usr/share/gvfs/mounts/ftps.mount
[Mount]
Type=ftps
Exec=/usr/lib/gvfs/gvfsd-ftp
AutoMount=false
Scheme=ftps
DefaultPort=21
HostnameIsInetAddress=true
