# squashfs-root/usr/share/gvfs/mounts/google.mount
[Mount]
Type=google-drive
Exec=/usr/lib/gvfs/gvfsd-google
AutoMount=false
