# squashfs-root/usr/share/gvfs/mounts/gphoto2.mount
[Mount]
Type=gphoto2
Exec=/usr/lib/gvfs/gvfsd-gphoto2
AutoMount=false
