# squashfs-root/usr/share/gvfs/mounts/http.mount
[Mount]
Type=http
Exec=/usr/lib/gvfs/gvfsd-http
AutoMount=true
DBusName=org.gtk.vfs.mountpoint_http
