# squashfs-root/usr/share/gvfs/mounts/localtest.mount
[Mount]
Type=localtest
Exec=/usr/lib/gvfs/gvfsd-localtest
AutoMount=false
