# squashfs-root/usr/share/gvfs/mounts/mtp.mount
[Mount]
Type=mtp
Exec=/usr/lib/gvfs/gvfsd-mtp
AutoMount=false
