# squashfs-root/usr/share/gvfs/mounts/recent.mount
[Mount]
Type=recent
Exec=/usr/lib/gvfs/gvfsd-recent
AutoMount=true

