# squashfs-root/usr/share/gvfs/mounts/sftp.mount
[Mount]
Type=sftp
Exec=/usr/lib/gvfs/gvfsd-sftp
AutoMount=false
Scheme=sftp
SchemeAliases=ssh
DefaultPort=22
HostnameIsInetAddress=true
