# squashfs-root/usr/share/gvfs/mounts/smb.mount
[Mount]
Type=smb-share
Exec=/usr/lib/gvfs/gvfsd-smb
AutoMount=false
Scheme=smb
