# squashfs-root/lib/systemd/system/apport-forward@.service
[Unit]
Description=Apport crash forwarding receiver
Requires=apport-forward.socket

[Service]
Type=oneshot
ExecStart=/usr/share/apport/apport
