# squashfs-root/lib/systemd/system/apport-autoreport.path
[Unit]
Description=Process error reports when automatic reporting is enabled (file watch)
ConditionPathExists=/var/lib/apport/autoreport

[Path]
PathExistsGlob=/var/crash/*.crash

[Install]
WantedBy=paths.target
