# squashfs-root/lib/systemd/system/apport-forward.socket
[Unit]
Description=Unix socket for apport crash forwarding
ConditionVirtualization=container

[Socket]
ListenStream=/run/apport.socket
SocketMode=0600
Accept=yes
MaxConnections=10
Backlog=5
PassCredentials=true

[Install]
WantedBy=sockets.target
