# squashfs-root/lib/systemd/system/friendly-recovery.target
[Unit]
Description=Friendly Recovery Mode
DefaultDependencies=no
Wants=friendly-recovery.service
After=friendly-recovery.service
Conflicts=shutdown.target
AllowIsolate=yes
