# squashfs-root/lib/systemd/system/kerneloops.service
[Unit]
Description=Tool to automatically collect and submit kernel crash signatures
After=remote-fs.target nss-lookup.target network-online.target time-sync.target

[Service]
EnvironmentFile=-/etc/default/kerneloops
Type=forking
ExecStartPre=/usr/sbin/kerneloops --test
ExecStart=/usr/sbin/kerneloops $DAEMON_ARGS
User=kernoops
Group=adm

[Install]
WantedBy=multi-user.target
