# squashfs-root/lib/systemd/system/spice-vdagentd.target
[Unit]
Description=Agent daemon for Spice guests
