# squashfs-root/lib/systemd/system/ureadahead-stop.timer
[Unit]
Description=Stop ureadahead data collection 45s after completed startup
DefaultDependencies=false
After=default.target
Conflicts=shutdown.target
Before=shutdown.target
ConditionVirtualization=no

[Timer]
OnActiveSec=45s
AccuracySec=1s
