# squashfs-root/lib/systemd/system/uuidd.socket
[Unit]
Description=UUID daemon activation socket

[Socket]
ListenStream=/run/uuidd/request

[Install]
WantedBy=sockets.target
