# squashfs-root/snap/gnome-3-28-1804/67/etc/xdg/autostart/at-spi-dbus-bus.desktop
[Desktop Entry]
Type=Application
Name=AT-SPI D-Bus Bus
Exec=/usr/lib/at-spi2-core/at-spi-bus-launcher --launch-immediately
NoDisplay=true
X-GNOME-AutoRestart=true
X-GNOME-Autostart-Phase=Initialization
