# squashfs-root/etc/rc.local
#!/bin/sh
sudo -i -u extern dbus-run-session -- gio set /home/extern/Desktop/ubiquity.desktop metadata::trusted yes
exit 0
