# squashfs-root/etc/systemd/system/snap-gtk\x2dcommon\x2dthemes-1313.mount
[Unit]
Description=Mount unit for gtk-common-themes, revision 1313
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/gtk-common-themes_1313.snap
Where=/snap/gtk-common-themes/1313
Type=squashfs
Options=nodev,ro,x-gdu.hide

[Install]
WantedBy=multi-user.target
