# squashfs-root/etc/systemd/system/snap-core-7270.mount
[Unit]
Description=Mount unit for core, revision 7270
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/core_7270.snap
Where=/snap/core/7270
Type=squashfs
Options=nodev,ro,x-gdu.hide

[Install]
WantedBy=multi-user.target
