# squashfs-root/etc/systemd/system/snap-gnome\x2dcharacters-296.mount
[Unit]
Description=Mount unit for gnome-characters, revision 296
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/gnome-characters_296.snap
Where=/snap/gnome-characters/296
Type=squashfs
Options=nodev,ro,x-gdu.hide

[Install]
WantedBy=multi-user.target
