# squashfs-root/var/lib/dpkg/info/accountsservice.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.Accounts.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/accounts-daemon.service
/usr
/usr/lib
/usr/lib/accountsservice
/usr/lib/accountsservice/accounts-daemon
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.Accounts.User.xml
/usr/share/dbus-1/interfaces/org.freedesktop.Accounts.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.Accounts.service
/usr/share/doc
/usr/share/doc/accountsservice
/usr/share/doc/accountsservice/README
/usr/share/doc/accountsservice/TODO
/usr/share/doc/accountsservice/copyright
/usr/share/language-tools
/usr/share/language-tools/language-options
/usr/share/language-tools/language-validate
/usr/share/language-tools/language2locale
/usr/share/language-tools/locale2papersize
/usr/share/language-tools/main-countries
/usr/share/language-tools/save-to-pam-env
/usr/share/language-tools/set-language-helper
/usr/share/language-tools/update-langlist
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.accounts.policy
/var
/var/lib
/var/lib/AccountsService
/var/lib/AccountsService/icons
/var/lib/AccountsService/users
/usr/share/doc/accountsservice/changelog.Debian.gz
