# squashfs-root/var/lib/dpkg/info/accountsservice.postrm
#!/bin/sh
set -e
# Automatically added by dh_systemd_start/11ubuntu1
if [ -d /run/systemd/system ]; then
	systemctl --system daemon-reload >/dev/null || true
fi
# End automatically added section
# Automatically added by dh_systemd_enable/11ubuntu1
if [ "$1" = "remove" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper mask 'accounts-daemon.service' >/dev/null
	fi
fi

if [ "$1" = "purge" ]; then
	if [ -x "/usr/bin/deb-systemd-helper" ]; then
		deb-systemd-helper purge 'accounts-daemon.service' >/dev/null
		deb-systemd-helper unmask 'accounts-daemon.service' >/dev/null
	fi
fi
# End automatically added section
