# squashfs-root/var/lib/dpkg/info/acpi-support.postinst
#!/bin/bash

case "$1" in
        configure)
                invoke-rc.d acpid restart>/dev/null
                ;;
        *)
                ;;
esac


