# squashfs-root/var/lib/dpkg/info/anacron.preinst
#!/bin/sh

set -e

if dpkg --compare-versions "$2" lt-nl 2.3-24; then
	deb-systemd-helper purge anacron-resume.service >/dev/null
	deb-systemd-helper unmask anacron-resume.service >/dev/null
fi


