# squashfs-root/var/lib/dpkg/info/anacron.list
/.
/etc
/etc/anacrontab
/etc/apm
/etc/apm/event.d
/etc/apm/event.d/anacron
/etc/cron.d
/etc/cron.d/anacron
/etc/cron.daily
/etc/cron.daily/0anacron
/etc/cron.monthly
/etc/cron.monthly/0anacron
/etc/cron.weekly
/etc/cron.weekly/0anacron
/etc/default
/etc/default/anacron
/etc/init
/etc/init/anacron.conf
/etc/init.d
/etc/init.d/anacron
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/anacron.service
/lib/systemd/system/anacron.timer
/usr
/usr/lib
/usr/lib/pm-utils
/usr/lib/pm-utils/power.d
/usr/lib/pm-utils/power.d/anacron
/usr/lib/pm-utils/sleep.d
/usr/lib/pm-utils/sleep.d/95anacron
/usr/sbin
/usr/sbin/anacron
/usr/share
/usr/share/doc
/usr/share/doc/anacron
/usr/share/doc/anacron/README.Debian
/usr/share/doc/anacron/README.gz
/usr/share/doc/anacron/TODO
/usr/share/doc/anacron/changelog.Debian.gz
/usr/share/doc/anacron/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/anacrontab.5.gz
/usr/share/man/man8
/usr/share/man/man8/anacron.8.gz
/var
/var/spool
/var/spool/anacron
