# squashfs-root/var/lib/dpkg/info/anacron.conffiles
/etc/anacrontab
/etc/apm/event.d/anacron
/etc/cron.d/anacron
/etc/cron.daily/0anacron
/etc/cron.monthly/0anacron
/etc/cron.weekly/0anacron
/etc/default/anacron
/etc/init.d/anacron
/etc/init/anacron.conf
