# squashfs-root/var/lib/dpkg/info/apport.list
/.
/etc
/etc/apport
/etc/apport/blacklist.d
/etc/apport/blacklist.d/README.blacklist
/etc/apport/blacklist.d/apport
/etc/apport/crashdb.conf
/etc/bash_completion.d
/etc/bash_completion.d/apport_completion
/etc/cron.daily
/etc/cron.daily/apport
/etc/default
/etc/default/apport
/etc/init.d
/etc/init.d/apport
/etc/logrotate.d
/etc/logrotate.d/apport
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/apport-autoreport.path
/lib/systemd/system/apport-autoreport.service
/lib/systemd/system/apport-forward.socket
/lib/systemd/system/apport-forward@.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/50-apport.rules
/usr
/usr/bin
/usr/bin/apport-bug
/usr/bin/apport-cli
/usr/bin/apport-unpack
/usr/lib
/usr/lib/pm-utils
/usr/lib/pm-utils/sleep.d
/usr/lib/pm-utils/sleep.d/000record-status
/usr/share
/usr/share/apport
/usr/share/apport/apport
/usr/share/apport/apport-checkreports
/usr/share/apport/apport.jar
/usr/share/apport/apportcheckresume
/usr/share/apport/dump_acpi_tables.py
/usr/share/apport/gcc_ice_hook
/usr/share/apport/general-hooks
/usr/share/apport/general-hooks/clickinfo.py
/usr/share/apport/general-hooks/cloud_archive.py
/usr/share/apport/general-hooks/generic.py
/usr/share/apport/general-hooks/parse_segv.py
/usr/share/apport/general-hooks/powerpc.py
/usr/share/apport/general-hooks/ubuntu-gnome.py
/usr/share/apport/general-hooks/ubuntu.py
/usr/share/apport/general-hooks/wayland_session.py
/usr/share/apport/is-enabled
/usr/share/apport/iwlwifi_error_dump
/usr/share/apport/java_uncaught_exception
/usr/share/apport/kernel_crashdump
/usr/share/apport/kernel_oops
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_apport.py
/usr/share/apport/package-hooks/source_debian-installer.py
/usr/share/apport/package-hooks/source_linux-nexus7.py
/usr/share/apport/package-hooks/source_linux.py
/usr/share/apport/package-hooks/source_ubiquity.py
/usr/share/apport/package_hook
/usr/share/apport/recoverable_problem
/usr/share/apport/root_info_wrapper
/usr/share/apport/testsuite
/usr/share/apport/testsuite/crash.class
/usr/share/apport/testsuite/crash.jar
/usr/share/apport/unkillable_shutdown
/usr/share/apport/whoopsie-upload-all
/usr/share/doc
/usr/share/doc/apport
/usr/share/doc/apport/NEWS.gz
/usr/share/doc/apport/README
/usr/share/doc/apport/copyright
/usr/share/doc/apport/crashdb-conf.txt.gz
/usr/share/doc/apport/package-hooks.txt.gz
/usr/share/doc/apport/symptoms.txt
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/32x32
/usr/share/icons/hicolor/32x32/apps
/usr/share/icons/hicolor/32x32/apps/apport.png
/usr/share/icons/hicolor/32x32/mimetypes
/usr/share/icons/hicolor/48x48
/usr/share/icons/hicolor/48x48/apps
/usr/share/icons/hicolor/48x48/apps/apport.png
/usr/share/icons/hicolor/48x48/mimetypes
/usr/share/icons/hicolor/64x64
/usr/share/icons/hicolor/64x64/apps
/usr/share/icons/hicolor/64x64/apps/apport.png
/usr/share/icons/hicolor/64x64/mimetypes
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/apport.svg
/usr/share/icons/hicolor/scalable/mimetypes
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/apport-bug.1.gz
/usr/share/man/man1/apport-cli.1.gz
/usr/share/man/man1/apport-unpack.1.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/com.ubuntu.apport.policy
/usr/share/python
/usr/share/python/runtime.d
/usr/share/python/runtime.d/apport.rtupdate
/usr/bin/apport-collect
/usr/bin/ubuntu-bug
/usr/share/apport/package-hooks/source_linux-meta.py
/usr/share/apport/package-hooks/source_linux-signed.py
/usr/share/doc/apport/changelog.Debian.gz
/usr/share/icons/hicolor/32x32/mimetypes/text-x-apport.png
/usr/share/icons/hicolor/48x48/mimetypes/text-x-apport.png
/usr/share/icons/hicolor/64x64/mimetypes/text-x-apport.png
/usr/share/icons/hicolor/scalable/mimetypes/text-x-apport.svg
/usr/share/man/man1/apport-collect.1.gz
/usr/share/man/man1/ubuntu-bug.1.gz
