# squashfs-root/var/lib/dpkg/info/apport.conffiles
/etc/apport/blacklist.d/README.blacklist
/etc/apport/blacklist.d/apport
/etc/apport/crashdb.conf
/etc/bash_completion.d/apport_completion
/etc/cron.daily/apport
/etc/default/apport
/etc/init.d/apport
/etc/logrotate.d/apport
