# squashfs-root/var/lib/dpkg/info/apt-xapian-index.list
/etc
/etc/cron.weekly
/etc/cron.weekly/apt-xapian-index
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.debian.AptXapianIndex.conf
