# squashfs-root/var/lib/dpkg/info/bluez.conffiles
/etc/bluetooth/input.conf
/etc/bluetooth/main.conf
/etc/bluetooth/network.conf
/etc/dbus-1/system.d/bluetooth.conf
/etc/init.d/bluetooth
