# squashfs-root/var/lib/dpkg/info/bolt.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.bolt.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/bolt.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/90-bolt.rules
/usr
/usr/bin
/usr/bin/boltctl
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/boltd
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.bolt.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.bolt.service
/usr/share/doc
/usr/share/doc/bolt
/usr/share/doc/bolt/changelog.Debian.gz
/usr/share/doc/bolt/copyright
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.bolt.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.bolt.rules
/var
/var/lib
/var/lib/boltd
