# squashfs-root/var/lib/dpkg/info/colord.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.ColorManager.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/colord.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/69-cd-sensors.rules
/lib/udev/rules.d/95-cd-devices.rules
/usr
/usr/bin
/usr/bin/cd-create-profile
/usr/bin/cd-fix-profile
/usr/bin/cd-iccdump
/usr/bin/cd-it8
/usr/bin/colormgr
/usr/lib
/usr/lib/colord
/usr/lib/colord/colord
/usr/lib/colord/colord-sane
/usr/lib/colord/colord-session
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/colord-session.service
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/colord.conf
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/colord-plugins
/usr/lib/x86_64-linux-gnu/colord-plugins/libcd_plugin_camera.so
/usr/lib/x86_64-linux-gnu/colord-plugins/libcd_plugin_sane.so
/usr/lib/x86_64-linux-gnu/colord-plugins/libcd_plugin_scanner.so
/usr/lib/x86_64-linux-gnu/colord-sensors
/usr/lib/x86_64-linux-gnu/colord-sensors/libcolord_sensor_colorhug.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libcolord_sensor_dtp94.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libcolord_sensor_dummy.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libcolord_sensor_huey.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libcolord_sensor_spark.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libdtp94-private.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libhuey-private.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libmunki-private.so
/usr/lib/x86_64-linux-gnu/colord-sensors/libospark-private.so
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/colormgr
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.ColorHelper.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Device.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Profile.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.Sensor.xml
/usr/share/dbus-1/interfaces/org.freedesktop.ColorManager.xml
/usr/share/dbus-1/services
/usr/share/dbus-1/services/org.freedesktop.ColorHelper.service
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.ColorManager.service
/usr/share/doc
/usr/share/doc/colord
/usr/share/doc/colord/copyright
/usr/share/glib-2.0
/usr/share/glib-2.0/schemas
/usr/share/glib-2.0/schemas/org.freedesktop.ColorHelper.gschema.xml
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/cd-create-profile.1.gz
/usr/share/man/man1/cd-fix-profile.1.gz
/usr/share/man/man1/cd-it8.1.gz
/usr/share/man/man1/colormgr.1.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.color.policy
/usr/share/doc/colord/changelog.Debian.gz
