# squashfs-root/var/lib/dpkg/info/cracklib-runtime.conffiles
/etc/cracklib/cracklib.conf
/etc/cron.daily/cracklib-runtime
/etc/logcheck/ignore.d.paranoid/cracklib-runtime
