# squashfs-root/var/lib/dpkg/info/cracklib-runtime.list
/.
/etc
/etc/cracklib
/etc/cracklib/cracklib.conf
/etc/cron.daily
/etc/cron.daily/cracklib-runtime
/etc/logcheck
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/cracklib-runtime
/usr
/usr/sbin
/usr/sbin/cracklib-check
/usr/sbin/cracklib-format
/usr/sbin/cracklib-packer
/usr/sbin/cracklib-unpacker
/usr/sbin/create-cracklib-dict
/usr/sbin/update-cracklib
/usr/share
/usr/share/dict
/usr/share/dict/cracklib-small
/usr/share/doc
/usr/share/doc/cracklib-runtime
/usr/share/doc/cracklib-runtime/copyright
/usr/share/doc/cracklib-runtime/cracklib-runtime.html
/usr/share/doc-base
/usr/share/doc-base/cracklib-runtime
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/cracklib-check.8.gz
/usr/share/man/man8/cracklib-format.8.gz
/usr/share/man/man8/create-cracklib-dict.8.gz
/usr/share/man/man8/update-cracklib.8.gz
/var
/var/cache
/var/cache/cracklib
/usr/share/doc/cracklib-runtime/changelog.Debian.gz
/usr/share/man/man8/cracklib-packer.8.gz
/usr/share/man/man8/cracklib-unpacker.8.gz
