# squashfs-root/var/lib/dpkg/info/cryptsetup.list
/.
/etc
/etc/cryptsetup-initramfs
/etc/cryptsetup-initramfs/conf-hook
/etc/default
/etc/default/cryptdisks
/etc/init.d
/etc/init.d/cryptdisks
/etc/init.d/cryptdisks-early
/lib
/lib/cryptsetup
/lib/cryptsetup/askpass
/lib/cryptsetup/checks
/lib/cryptsetup/checks/blkid
/lib/cryptsetup/checks/ext2
/lib/cryptsetup/checks/swap
/lib/cryptsetup/checks/un_blkid
/lib/cryptsetup/checks/xfs
/lib/cryptsetup/cryptdisks.functions
/lib/cryptsetup/scripts
/lib/cryptsetup/scripts/decrypt_derived
/lib/cryptsetup/scripts/decrypt_gnupg
/lib/cryptsetup/scripts/decrypt_keyctl
/lib/cryptsetup/scripts/decrypt_openct
/lib/cryptsetup/scripts/decrypt_opensc
/lib/cryptsetup/scripts/decrypt_ssl
/lib/cryptsetup/scripts/passdev
/sbin
/sbin/cryptdisks_start
/sbin/cryptdisks_stop
/usr
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/cryptsetup.py
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/cryptdisks
/usr/share/bug
/usr/share/bug/cryptsetup
/usr/share/cryptsetup
/usr/share/cryptsetup/initramfs
/usr/share/cryptsetup/initramfs/bin
/usr/share/cryptsetup/initramfs/bin/cryptroot-unlock
/usr/share/doc
/usr/share/doc/cryptsetup
/usr/share/doc/cryptsetup/AUTHORS
/usr/share/doc/cryptsetup/FAQ.gz
/usr/share/doc/cryptsetup/README.Debian.gz
/usr/share/doc/cryptsetup/README.gnupg
/usr/share/doc/cryptsetup/README.initramfs.gz
/usr/share/doc/cryptsetup/README.keyctl
/usr/share/doc/cryptsetup/README.openct.gz
/usr/share/doc/cryptsetup/README.opensc.gz
/usr/share/doc/cryptsetup/TODO.Debian
/usr/share/doc/cryptsetup/copyright
/usr/share/doc/cryptsetup/examples
/usr/share/doc/cryptsetup/examples/gen-ssl-key
/usr/share/doc/cryptsetup/v1.0.7-ReleaseNotes
/usr/share/doc/cryptsetup/v1.1.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.1.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.1.2-ReleaseNotes
/usr/share/doc/cryptsetup/v1.1.3-ReleaseNotes
/usr/share/doc/cryptsetup/v1.2.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.3.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.3.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.4.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.4.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.4.2-ReleaseNotes
/usr/share/doc/cryptsetup/v1.4.3-ReleaseNotes
/usr/share/doc/cryptsetup/v1.5.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.5.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v1.6.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.2-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.3-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.4-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.5-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.6-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.7-ReleaseNotes
/usr/share/doc/cryptsetup/v1.6.8-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.0-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.1-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.2-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.3-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.4-ReleaseNotes
/usr/share/doc/cryptsetup/v1.7.5-ReleaseNotes
/usr/share/doc/cryptsetup/v2.0.0-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v2.0.1-ReleaseNotes.gz
/usr/share/doc/cryptsetup/v2.0.2-ReleaseNotes.gz
/usr/share/initramfs-tools
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/initramfs-tools/conf-hooks.d/cryptsetup
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/cryptgnupg
/usr/share/initramfs-tools/hooks/cryptkeyctl
/usr/share/initramfs-tools/hooks/cryptopenct
/usr/share/initramfs-tools/hooks/cryptopensc
/usr/share/initramfs-tools/hooks/cryptpassdev
/usr/share/initramfs-tools/hooks/cryptroot
/usr/share/initramfs-tools/hooks/cryptroot-unlock
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/local-block
/usr/share/initramfs-tools/scripts/local-block/cryptroot
/usr/share/initramfs-tools/scripts/local-bottom
/usr/share/initramfs-tools/scripts/local-bottom/cryptopensc
/usr/share/initramfs-tools/scripts/local-top
/usr/share/initramfs-tools/scripts/local-top/cryptopensc
/usr/share/initramfs-tools/scripts/local-top/cryptroot
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cryptsetup
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/crypttab.5.gz
/usr/share/man/man8
/usr/share/man/man8/cryptdisks_start.8.gz
/usr/share/man/man8/cryptdisks_stop.8.gz
/usr/share/doc/cryptsetup/NEWS.Debian.gz
/usr/share/doc/cryptsetup/changelog.Debian.gz
