# squashfs-root/var/lib/dpkg/info/cups-browsed.conffiles
/etc/apparmor.d/usr.sbin.cups-browsed
/etc/cups/cups-browsed.conf
/etc/init.d/cups-browsed
