# squashfs-root/var/lib/dpkg/info/flatpak.list
/.
/etc
/etc/X11
/etc/X11/Xsession.d
/etc/X11/Xsession.d/20flatpak
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.Flatpak.SystemHelper.conf
/etc/profile.d
/etc/profile.d/flatpak.sh
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/flatpak-system-helper.service
/usr
/usr/bin
/usr/bin/flatpak
/usr/lib
/usr/lib/flatpak
/usr/lib/flatpak/flatpak-dbus-proxy
/usr/lib/flatpak/flatpak-portal
/usr/lib/flatpak/flatpak-session-helper
/usr/lib/flatpak/flatpak-system-helper
/usr/lib/systemd
/usr/lib/systemd/user
/usr/lib/systemd/user/dbus.service.d
/usr/lib/systemd/user/dbus.service.d/flatpak.conf
/usr/lib/systemd/user/flatpak-portal.service
/usr/lib/systemd/user/flatpak-session-helper.service
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/flatpak
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/services/org.freedesktop.Flatpak.service
/usr/share/dbus-1/services/org.freedesktop.portal.Flatpak.service
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.Flatpak.SystemHelper.service
/usr/share/doc
/usr/share/doc/flatpak
/usr/share/doc/flatpak/NEWS.gz
/usr/share/doc/flatpak/README.md
/usr/share/doc/flatpak/changelog.Debian.gz
/usr/share/doc/flatpak/copyright
/usr/share/doc/flatpak/docbook.css
/usr/share/doc/flatpak/examples
/usr/share/doc/flatpak/examples/flatpak-bisect.gz
/usr/share/doc/flatpak/examples/flatpak-coredumpctl
/usr/share/doc/flatpak/flatpak-docs.html
/usr/share/doc-base
/usr/share/doc-base/flatpak-commands
/usr/share/flatpak
/usr/share/flatpak/triggers
/usr/share/flatpak/triggers/desktop-database.trigger
/usr/share/flatpak/triggers/gtk-icon-cache.trigger
/usr/share/flatpak/triggers/mime-database.trigger
/usr/share/gdm
/usr/share/gdm/env.d
/usr/share/gdm/env.d/flatpak.env
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/flatpak
/usr/share/locale
/usr/share/locale/cs
/usr/share/locale/cs/LC_MESSAGES
/usr/share/locale/cs/LC_MESSAGES/flatpak.mo
/usr/share/locale/de
/usr/share/locale/de/LC_MESSAGES
/usr/share/locale/de/LC_MESSAGES/flatpak.mo
/usr/share/locale/es
/usr/share/locale/es/LC_MESSAGES
/usr/share/locale/es/LC_MESSAGES/flatpak.mo
/usr/share/locale/gl
/usr/share/locale/gl/LC_MESSAGES
/usr/share/locale/gl/LC_MESSAGES/flatpak.mo
/usr/share/locale/hu
/usr/share/locale/hu/LC_MESSAGES
/usr/share/locale/hu/LC_MESSAGES/flatpak.mo
/usr/share/locale/id
/usr/share/locale/id/LC_MESSAGES
/usr/share/locale/id/LC_MESSAGES/flatpak.mo
/usr/share/locale/pl
/usr/share/locale/pl/LC_MESSAGES
/usr/share/locale/pl/LC_MESSAGES/flatpak.mo
/usr/share/locale/pt_BR
/usr/share/locale/pt_BR/LC_MESSAGES
/usr/share/locale/pt_BR/LC_MESSAGES/flatpak.mo
/usr/share/locale/ru
/usr/share/locale/ru/LC_MESSAGES
/usr/share/locale/ru/LC_MESSAGES/flatpak.mo
/usr/share/locale/sk
/usr/share/locale/sk/LC_MESSAGES
/usr/share/locale/sk/LC_MESSAGES/flatpak.mo
/usr/share/locale/sv
/usr/share/locale/sv/LC_MESSAGES
/usr/share/locale/sv/LC_MESSAGES/flatpak.mo
/usr/share/locale/tr
/usr/share/locale/tr/LC_MESSAGES
/usr/share/locale/tr/LC_MESSAGES/flatpak.mo
/usr/share/locale/uk
/usr/share/locale/uk/LC_MESSAGES
/usr/share/locale/uk/LC_MESSAGES/flatpak.mo
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/flatpak-build-bundle.1.gz
/usr/share/man/man1/flatpak-build-commit-from.1.gz
/usr/share/man/man1/flatpak-build-export.1.gz
/usr/share/man/man1/flatpak-build-finish.1.gz
/usr/share/man/man1/flatpak-build-import-bundle.1.gz
/usr/share/man/man1/flatpak-build-init.1.gz
/usr/share/man/man1/flatpak-build-sign.1.gz
/usr/share/man/man1/flatpak-build-update-repo.1.gz
/usr/share/man/man1/flatpak-build.1.gz
/usr/share/man/man1/flatpak-config.1.gz
/usr/share/man/man1/flatpak-create-usb.1.gz
/usr/share/man/man1/flatpak-document-export.1.gz
/usr/share/man/man1/flatpak-document-info.1.gz
/usr/share/man/man1/flatpak-document-list.1.gz
/usr/share/man/man1/flatpak-document-unexport.1.gz
/usr/share/man/man1/flatpak-enter.1.gz
/usr/share/man/man1/flatpak-info.1.gz
/usr/share/man/man1/flatpak-install.1.gz
/usr/share/man/man1/flatpak-list.1.gz
/usr/share/man/man1/flatpak-make-current.1.gz
/usr/share/man/man1/flatpak-override.1.gz
/usr/share/man/man1/flatpak-permission-list.1.gz
/usr/share/man/man1/flatpak-permission-remove.1.gz
/usr/share/man/man1/flatpak-permission-reset.1.gz
/usr/share/man/man1/flatpak-permission-show.1.gz
/usr/share/man/man1/flatpak-ps.1.gz
/usr/share/man/man1/flatpak-remote-add.1.gz
/usr/share/man/man1/flatpak-remote-delete.1.gz
/usr/share/man/man1/flatpak-remote-info.1.gz
/usr/share/man/man1/flatpak-remote-ls.1.gz
/usr/share/man/man1/flatpak-remote-modify.1.gz
/usr/share/man/man1/flatpak-remotes.1.gz
/usr/share/man/man1/flatpak-repair.1.gz
/usr/share/man/man1/flatpak-repo.1.gz
/usr/share/man/man1/flatpak-run.1.gz
/usr/share/man/man1/flatpak-search.1.gz
/usr/share/man/man1/flatpak-uninstall.1.gz
/usr/share/man/man1/flatpak-update.1.gz
/usr/share/man/man1/flatpak.1.gz
/usr/share/man/man5
/usr/share/man/man5/flatpak-flatpakref.5.gz
/usr/share/man/man5/flatpak-flatpakrepo.5.gz
/usr/share/man/man5/flatpak-installation.5.gz
/usr/share/man/man5/flatpak-metadata.5.gz
/usr/share/man/man5/flatpak-remote.5.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.Flatpak.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.Flatpak.rules
/usr/share/zsh
/usr/share/zsh/vendor-completions
/usr/share/zsh/vendor-completions/_flatpak
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/org.freedesktop.Flatpak.pkla
