# squashfs-root/var/lib/dpkg/info/friendly-recovery.list
/.
/lib
/lib/recovery-mode
/lib/recovery-mode/l10n.sh
/lib/recovery-mode/options
/lib/recovery-mode/options/apt-snapshots
/lib/recovery-mode/options/clean
/lib/recovery-mode/options/dpkg
/lib/recovery-mode/options/failsafeX
/lib/recovery-mode/options/fsck
/lib/recovery-mode/options/grub
/lib/recovery-mode/options/network
/lib/recovery-mode/options/root
/lib/recovery-mode/options/system-summary
/lib/recovery-mode/recovery-menu
/lib/systemd
/lib/systemd/system
/lib/systemd/system/friendly-recovery.service
/lib/systemd/system/friendly-recovery.target
/lib/systemd/system-generators
/lib/systemd/system-generators/friendly-recovery
/usr
/usr/share
/usr/share/doc
/usr/share/doc/friendly-recovery
/usr/share/doc/friendly-recovery/README
/usr/share/doc/friendly-recovery/changelog.gz
/usr/share/doc/friendly-recovery/copyright
