# squashfs-root/var/lib/dpkg/info/friendly-recovery.preinst
#!/bin/sh
set -e

# it is now a static service
systemctl disable friendly-recovery.service || :

# Automatically added by dh_installdeb/11.1.6ubuntu2
dpkg-maintscript-helper rm_conffile /etc/init/friendly-recovery.conf 0.2.36\~ -- "$@"
# End automatically added section


exit 0
