# squashfs-root/var/lib/dpkg/info/fwupdate.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/fwupdate-cleanup.service
/usr
/usr/bin
/usr/bin/fwupdate
/usr/lib
/usr/lib/fwupdate
/usr/lib/fwupdate/cleanup
/usr/lib/fwupdate/fwupx64.efi
/usr/lib/fwupdate/install
/usr/share
/usr/share/doc
/usr/share/doc/fwupdate
/usr/share/doc/fwupdate/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fwupdate
/usr/share/locale
/usr/share/locale/en
/usr/share/locale/en/fwupdate.po
/usr/share/locale/en/libfwup.po
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fwupdate.1.gz
/usr/share/man/man3
/usr/share/man/man3/libfwup.3.gz
/var
/var/lib
/var/lib/fwupdate
/usr/share/doc/fwupdate/changelog.Debian.gz
/usr/share/man/man3/fwup_clear_status.3.gz
/usr/share/man/man3/fwup_get_fw_type.3.gz
/usr/share/man/man3/fwup_get_fw_version.3.gz
/usr/share/man/man3/fwup_get_guid.3.gz
/usr/share/man/man3/fwup_get_last_attempt_info.3.gz
/usr/share/man/man3/fwup_get_lowest_supported_version.3.gz
/usr/share/man/man3/fwup_get_ux_capsule_info.3.gz
/usr/share/man/man3/fwup_resource_free.3.gz
/usr/share/man/man3/fwup_resource_iter_create.3.gz
/usr/share/man/man3/fwup_resource_iter_destroy.3.gz
/usr/share/man/man3/fwup_resource_iter_next.3.gz
/usr/share/man/man3/fwup_set_guid.3.gz
/usr/share/man/man3/fwup_set_guid_forced.3.gz
/usr/share/man/man3/fwup_set_up_update.3.gz
/usr/share/man/man3/fwup_supported.3.gz
/usr/share/man/man3/libfwup.h.3.gz
