# squashfs-root/var/lib/dpkg/info/gdm3.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/gdm.conf
/etc/gdm3
/etc/gdm3/Init
/etc/gdm3/Init/Default
/etc/gdm3/PostLogin
/etc/gdm3/PostLogin/Default.sample
/etc/gdm3/PostSession
/etc/gdm3/PostSession/Default
/etc/gdm3/PreSession
/etc/gdm3/PreSession/Default
/etc/gdm3/Prime
/etc/gdm3/Prime/Default
/etc/gdm3/PrimeOff
/etc/gdm3/PrimeOff/Default
/etc/gdm3/Xsession
/etc/gdm3/config-error-dialog.sh
/etc/gdm3/custom.conf
/etc/init.d
/etc/init.d/gdm3
/etc/insserv.conf.d
/etc/insserv.conf.d/gdm3
/etc/pam.d
/etc/pam.d/gdm-autologin
/etc/pam.d/gdm-fingerprint
/etc/pam.d/gdm-launch-environment
/etc/pam.d/gdm-password
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/gdm.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/61-gdm.rules
/lib/x86_64-linux-gnu
/lib/x86_64-linux-gnu/security
/lib/x86_64-linux-gnu/security/pam_gdm.so
/usr
/usr/bin
/usr/bin/gdm-screenshot
/usr/lib
/usr/lib/gdm3
/usr/lib/gdm3/gdm-disable-wayland
/usr/lib/gdm3/gdm-host-chooser
/usr/lib/gdm3/gdm-session-worker
/usr/lib/gdm3/gdm-simple-chooser
/usr/lib/gdm3/gdm-wayland-session
/usr/lib/gdm3/gdm-x-session
/usr/sbin
/usr/sbin/gdm3
/usr/share
/usr/share/dconf
/usr/share/dconf/profile
/usr/share/dconf/profile/gdm
/usr/share/doc
/usr/share/doc/gdm3
/usr/share/doc/gdm3/AUTHORS
/usr/share/doc/gdm3/NEWS.gz
/usr/share/doc/gdm3/README
/usr/share/doc/gdm3/README.Debian
/usr/share/doc/gdm3/copyright
/usr/share/gdm
/usr/share/gdm/BuiltInSessions
/usr/share/gdm/dconf
/usr/share/gdm/dconf/00-upstream-settings
/usr/share/gdm/dconf/locks
/usr/share/gdm/dconf/locks/00-upstream-settings-locks
/usr/share/gdm/default.pa
/usr/share/gdm/gdb-cmd
/usr/share/gdm/gdm.schemas
/usr/share/gdm/generate-config
/usr/share/gdm/greeter
/usr/share/gdm/greeter/applications
/usr/share/gdm/greeter/applications/mime-dummy-handler.desktop
/usr/share/gdm/greeter/applications/mimeapps.list
/usr/share/gdm/greeter/autostart
/usr/share/gdm/greeter/autostart/orca-autostart.desktop
/usr/share/gdm/greeter.dconf-defaults
/usr/share/gdm/greeter.dconf-defaults.md5sum
/usr/share/gdm/locale.alias
/usr/share/gnome-session
/usr/share/gnome-session/sessions
/usr/share/gnome-session/sessions/gnome-login.session
/usr/share/help
/usr/share/help/C
/usr/share/help/C/gdm
/usr/share/help/C/gdm/index.docbook
/usr/share/help/C/gdm/legal.xml
/usr/share/help/ca
/usr/share/help/ca/gdm
/usr/share/help/cs
/usr/share/help/cs/gdm
/usr/share/help/de
/usr/share/help/de/gdm
/usr/share/help/el
/usr/share/help/el/gdm
/usr/share/help/en_GB
/usr/share/help/en_GB/gdm
/usr/share/help/es
/usr/share/help/es/gdm
/usr/share/help/fr
/usr/share/help/fr/gdm
/usr/share/help/gl
/usr/share/help/gl/gdm
/usr/share/help/hu
/usr/share/help/hu/gdm
/usr/share/help/id
/usr/share/help/id/gdm
/usr/share/help/it
/usr/share/help/it/gdm
/usr/share/help/ko
/usr/share/help/ko/gdm
/usr/share/help/oc
/usr/share/help/oc/gdm
/usr/share/help/pt_BR
/usr/share/help/pt_BR/gdm
/usr/share/help/ro
/usr/share/help/ro/gdm
/usr/share/help/ru
/usr/share/help/ru/gdm
/usr/share/help/sl
/usr/share/help/sl/gdm
/usr/share/help/sv
/usr/share/help/sv/gdm
/usr/share/help/te
/usr/share/help/te/gdm
/usr/share/help/tr
/usr/share/help/tr/gdm
/usr/share/help/uk
/usr/share/help/uk/gdm
/usr/share/help/zh_CN
/usr/share/help/zh_CN/gdm
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/16x16
/usr/share/icons/hicolor/16x16/apps
/usr/share/icons/hicolor/16x16/apps/gdm-xnest.png
/usr/share/icons/hicolor/32x32
/usr/share/icons/hicolor/32x32/apps
/usr/share/icons/hicolor/32x32/apps/gdm-setup.png
/usr/share/icons/hicolor/32x32/apps/gdm-xnest.png
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/gdm3
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/gdm-screenshot.1.gz
/usr/share/man/man8
/usr/share/man/man8/gdm3.8.gz
/usr/share/pixmaps
/usr/share/pixmaps/gdm-foot-logo.png
/usr/share/pixmaps/gdm-setup.png
/usr/share/pixmaps/gdm-xnest.png
/usr/share/pixmaps/gdm.png
/usr/share/pixmaps/nobody.png
/usr/share/pixmaps/nohost.png
/var
/var/cache
/var/cache/gdm
/var/lib
/var/lib/gdm3
/var/log
/var/log/gdm3
/lib/systemd/system/gdm3.service
/usr/share/doc/gdm3/changelog.Debian.gz
/usr/share/gdm/dconf/90-debian-settings
/usr/share/help/ca/gdm/index.docbook
/usr/share/help/ca/gdm/legal.xml
/usr/share/help/cs/gdm/index.docbook
/usr/share/help/cs/gdm/legal.xml
/usr/share/help/de/gdm/index.docbook
/usr/share/help/de/gdm/legal.xml
/usr/share/help/el/gdm/index.docbook
/usr/share/help/el/gdm/legal.xml
/usr/share/help/en_GB/gdm/index.docbook
/usr/share/help/en_GB/gdm/legal.xml
/usr/share/help/es/gdm/index.docbook
/usr/share/help/es/gdm/legal.xml
/usr/share/help/fr/gdm/index.docbook
/usr/share/help/fr/gdm/legal.xml
/usr/share/help/gl/gdm/index.docbook
/usr/share/help/gl/gdm/legal.xml
/usr/share/help/hu/gdm/index.docbook
/usr/share/help/hu/gdm/legal.xml
/usr/share/help/id/gdm/index.docbook
/usr/share/help/id/gdm/legal.xml
/usr/share/help/it/gdm/index.docbook
/usr/share/help/it/gdm/legal.xml
/usr/share/help/ko/gdm/index.docbook
/usr/share/help/ko/gdm/legal.xml
/usr/share/help/oc/gdm/index.docbook
/usr/share/help/oc/gdm/legal.xml
/usr/share/help/pt_BR/gdm/index.docbook
/usr/share/help/pt_BR/gdm/legal.xml
/usr/share/help/ro/gdm/index.docbook
/usr/share/help/ro/gdm/legal.xml
/usr/share/help/ru/gdm/index.docbook
/usr/share/help/ru/gdm/legal.xml
/usr/share/help/sl/gdm/index.docbook
/usr/share/help/sl/gdm/legal.xml
/usr/share/help/sv/gdm/index.docbook
/usr/share/help/sv/gdm/legal.xml
/usr/share/help/te/gdm/index.docbook
/usr/share/help/te/gdm/legal.xml
/usr/share/help/tr/gdm/index.docbook
/usr/share/help/tr/gdm/legal.xml
/usr/share/help/uk/gdm/index.docbook
/usr/share/help/uk/gdm/legal.xml
/usr/share/help/zh_CN/gdm/index.docbook
/usr/share/help/zh_CN/gdm/legal.xml
