# squashfs-root/var/lib/dpkg/info/gdm3.conffiles
/etc/dbus-1/system.d/gdm.conf
/etc/gdm3/Init/Default
/etc/gdm3/PostLogin/Default.sample
/etc/gdm3/PostSession/Default
/etc/gdm3/PreSession/Default
/etc/gdm3/Prime/Default
/etc/gdm3/PrimeOff/Default
/etc/gdm3/Xsession
/etc/gdm3/config-error-dialog.sh
/etc/gdm3/custom.conf
/etc/init.d/gdm3
/etc/insserv.conf.d/gdm3
/etc/pam.d/gdm-autologin
/etc/pam.d/gdm-fingerprint
/etc/pam.d/gdm-launch-environment
/etc/pam.d/gdm-password
