# squashfs-root/var/lib/dpkg/info/grub-common.list
/.
/etc
/etc/bash_completion.d
/etc/bash_completion.d/grub
/etc/grub.d
/etc/grub.d/00_header
/etc/grub.d/05_debian_theme
/etc/grub.d/10_linux
/etc/grub.d/20_linux_xen
/etc/grub.d/30_os-prober
/etc/grub.d/30_uefi-firmware
/etc/grub.d/40_custom
/etc/grub.d/41_custom
/etc/grub.d/README
/etc/init.d
/etc/init.d/grub-common
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_grub-common
/usr
/usr/bin
/usr/bin/grub-editenv
/usr/bin/grub-file
/usr/bin/grub-fstest
/usr/bin/grub-glue-efi
/usr/bin/grub-kbdcomp
/usr/bin/grub-menulst2cfg
/usr/bin/grub-mkfont
/usr/bin/grub-mkimage
/usr/bin/grub-mklayout
/usr/bin/grub-mknetdir
/usr/bin/grub-mkpasswd-pbkdf2
/usr/bin/grub-mkrelpath
/usr/bin/grub-mkrescue
/usr/bin/grub-mkstandalone
/usr/bin/grub-mount
/usr/bin/grub-render-label
/usr/bin/grub-script-check
/usr/bin/grub-syslinux2cfg
/usr/lib
/usr/lib/grub
/usr/sbin
/usr/sbin/grub-macbless
/usr/sbin/grub-mkconfig
/usr/sbin/grub-mkdevicemap
/usr/sbin/grub-probe
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_grub2.py
/usr/share/bug
/usr/share/bug/grub-common
/usr/share/bug/grub-common/presubj
/usr/share/bug/grub-common/script
/usr/share/doc
/usr/share/doc/grub-common
/usr/share/doc/grub-common/AUTHORS
/usr/share/doc/grub-common/NEWS.Debian.gz
/usr/share/doc/grub-common/NEWS.gz
/usr/share/doc/grub-common/README
/usr/share/doc/grub-common/THANKS
/usr/share/doc/grub-common/TODO
/usr/share/doc/grub-common/changelog.Debian.gz
/usr/share/doc/grub-common/copyright
/usr/share/doc/grub-common/examples
/usr/share/doc/grub-common/examples/grub.cfg
/usr/share/grub
/usr/share/grub/ascii.h
/usr/share/grub/ascii.pf2
/usr/share/grub/canonical-uefi-ca.crt
/usr/share/grub/euro.pf2
/usr/share/grub/grub-check-signatures
/usr/share/grub/grub-mkconfig_lib
/usr/share/grub/unicode.pf2
/usr/share/grub/widthspec.h
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/grub-editenv.1.gz
/usr/share/man/man1/grub-file.1.gz
/usr/share/man/man1/grub-fstest.1.gz
/usr/share/man/man1/grub-glue-efi.1.gz
/usr/share/man/man1/grub-kbdcomp.1.gz
/usr/share/man/man1/grub-menulst2cfg.1.gz
/usr/share/man/man1/grub-mkfont.1.gz
/usr/share/man/man1/grub-mkimage.1.gz
/usr/share/man/man1/grub-mklayout.1.gz
/usr/share/man/man1/grub-mknetdir.1.gz
/usr/share/man/man1/grub-mkpasswd-pbkdf2.1.gz
/usr/share/man/man1/grub-mkrelpath.1.gz
/usr/share/man/man1/grub-mkrescue.1.gz
/usr/share/man/man1/grub-mkstandalone.1.gz
/usr/share/man/man1/grub-mount.1.gz
/usr/share/man/man1/grub-render-label.1.gz
/usr/share/man/man1/grub-script-check.1.gz
/usr/share/man/man1/grub-syslinux2cfg.1.gz
/usr/share/man/man8
/usr/share/man/man8/grub-macbless.8.gz
/usr/share/man/man8/grub-mkconfig.8.gz
/usr/share/man/man8/grub-mkdevicemap.8.gz
/usr/share/man/man8/grub-probe.8.gz
/var
/var/lib
/var/lib/grub
/var/lib/grub/ucf
/usr/lib/grub/grub-mkconfig_lib
