# squashfs-root/var/lib/dpkg/info/iio-sensor-proxy.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/net.hadess.SensorProxy.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/iio-sensor-proxy.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/80-iio-sensor-proxy.rules
/usr
/usr/bin
/usr/bin/monitor-sensor
/usr/sbin
/usr/sbin/iio-sensor-proxy
/usr/share
/usr/share/doc
/usr/share/doc/iio-sensor-proxy
/usr/share/doc/iio-sensor-proxy/changelog.Debian.gz
/usr/share/doc/iio-sensor-proxy/copyright
/usr/share/gtk-doc
/usr/share/gtk-doc/html
/usr/share/gtk-doc/html/iio-sensor-proxy
/usr/share/gtk-doc/html/iio-sensor-proxy/gdbus-net.hadess.SensorProxy.Compass.html
/usr/share/gtk-doc/html/iio-sensor-proxy/gdbus-net.hadess.SensorProxy.html
/usr/share/gtk-doc/html/iio-sensor-proxy/home.png
/usr/share/gtk-doc/html/iio-sensor-proxy/iio-sensor-proxy.devhelp2
/usr/share/gtk-doc/html/iio-sensor-proxy/index.html
/usr/share/gtk-doc/html/iio-sensor-proxy/ix01.html
/usr/share/gtk-doc/html/iio-sensor-proxy/left-insensitive.png
/usr/share/gtk-doc/html/iio-sensor-proxy/left.png
/usr/share/gtk-doc/html/iio-sensor-proxy/license.html
/usr/share/gtk-doc/html/iio-sensor-proxy/ref-dbus.html
/usr/share/gtk-doc/html/iio-sensor-proxy/right-insensitive.png
/usr/share/gtk-doc/html/iio-sensor-proxy/right.png
/usr/share/gtk-doc/html/iio-sensor-proxy/style.css
/usr/share/gtk-doc/html/iio-sensor-proxy/up-insensitive.png
/usr/share/gtk-doc/html/iio-sensor-proxy/up.png
