# squashfs-root/var/lib/dpkg/info/irqbalance.conffiles
/etc/default/irqbalance
/etc/init.d/irqbalance
