# squashfs-root/var/lib/dpkg/info/irqbalance.list
/.
/etc
/etc/default
/etc/default/irqbalance
/etc/init.d
/etc/init.d/irqbalance
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/irqbalance.service
/usr
/usr/sbin
/usr/sbin/irqbalance
/usr/sbin/irqbalance-ui
/usr/share
/usr/share/doc
/usr/share/doc/irqbalance
/usr/share/doc/irqbalance/changelog.Debian.gz
/usr/share/doc/irqbalance/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/irqbalance
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/irqbalance.1.gz
