# squashfs-root/var/lib/dpkg/info/kerneloops.conffiles
/etc/dbus-1/system.d/kerneloops.conf
/etc/default/kerneloops
/etc/init.d/kerneloops
/etc/kerneloops.conf
