# squashfs-root/var/lib/dpkg/info/lightdm.conffiles
/etc/apparmor.d/abstractions/lightdm
/etc/apparmor.d/abstractions/lightdm_chromium-browser
/etc/apparmor.d/lightdm-guest-session
/etc/dbus-1/system.d/org.freedesktop.DisplayManager.conf
/etc/init.d/lightdm
/etc/init/lightdm.conf
/etc/lightdm/users.conf
/etc/logrotate.d/lightdm
/etc/pam.d/lightdm
/etc/pam.d/lightdm-autologin
/etc/pam.d/lightdm-greeter
