# squashfs-root/var/lib/dpkg/info/mlocate.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/mlocate
/etc/updatedb.conf
/usr
/usr/bin
/usr/bin/mlocate
/usr/bin/updatedb.mlocate
/usr/share
/usr/share/doc
/usr/share/doc/mlocate
/usr/share/doc/mlocate/AUTHORS
/usr/share/doc/mlocate/NEWS.gz
/usr/share/doc/mlocate/README
/usr/share/doc/mlocate/TODO.Debian
/usr/share/doc/mlocate/changelog.Debian.gz
/usr/share/doc/mlocate/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mlocate.1.gz
/usr/share/man/man5
/usr/share/man/man5/mlocate.db.5.gz
/usr/share/man/man5/updatedb.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/updatedb.mlocate.8.gz
/var
/var/lib
/var/lib/mlocate
