# squashfs-root/var/lib/dpkg/info/modemmanager.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.ModemManager1.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/ModemManager.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/77-mm-cinterion-port-types.rules
/lib/udev/rules.d/77-mm-dell-port-types.rules
/lib/udev/rules.d/77-mm-ericsson-mbm.rules
/lib/udev/rules.d/77-mm-fibocom-port-types.rules
/lib/udev/rules.d/77-mm-haier-port-types.rules
/lib/udev/rules.d/77-mm-huawei-net-port-types.rules
/lib/udev/rules.d/77-mm-longcheer-port-types.rules
/lib/udev/rules.d/77-mm-mtk-port-types.rules
/lib/udev/rules.d/77-mm-nokia-port-types.rules
/lib/udev/rules.d/77-mm-pcmcia-device-blacklist.rules
/lib/udev/rules.d/77-mm-qdl-device-blacklist.rules
/lib/udev/rules.d/77-mm-sierra.rules
/lib/udev/rules.d/77-mm-simtech-port-types.rules
/lib/udev/rules.d/77-mm-telit-port-types.rules
/lib/udev/rules.d/77-mm-ublox-port-types.rules
/lib/udev/rules.d/77-mm-usb-device-blacklist.rules
/lib/udev/rules.d/77-mm-usb-serial-adapters-greylist.rules
/lib/udev/rules.d/77-mm-x22x-port-types.rules
/lib/udev/rules.d/77-mm-zte-port-types.rules
/lib/udev/rules.d/80-mm-candidate.rules
/usr
/usr/bin
/usr/bin/mmcli
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/ModemManager
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-altair-lte.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-anydata.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-cinterion.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-dell.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-ericsson-mbm.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-fibocom.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-generic.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-haier.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-huawei.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-iridium.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-linktop.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-longcheer.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-motorola.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-mtk.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-nokia-icera.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-nokia.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-novatel-lte.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-novatel.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-option-hso.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-option.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-pantech.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-quectel.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-samsung.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-sierra-legacy.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-sierra.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-simtech.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-telit.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-thuraya.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-ublox.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-via.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-wavecom.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-x22x.so
/usr/lib/x86_64-linux-gnu/ModemManager/libmm-plugin-zte.so
/usr/sbin
/usr/sbin/ModemManager
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/mmcli
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.ModemManager1.service
/usr/share/doc
/usr/share/doc/modemmanager
/usr/share/doc/modemmanager/AUTHORS
/usr/share/doc/modemmanager/NEWS.gz
/usr/share/doc/modemmanager/README
/usr/share/doc/modemmanager/copyright
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/22x22
/usr/share/icons/hicolor/22x22/apps
/usr/share/icons/hicolor/22x22/apps/ModemManager.png
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mmcli.1.gz
/usr/share/man/man8
/usr/share/man/man8/ModemManager.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.ModemManager1.policy
/usr/share/doc/modemmanager/changelog.Debian.gz
