# squashfs-root/var/lib/dpkg/info/netplan.io.list
/.
/etc
/etc/netplan
/lib
/lib/netplan
/lib/netplan/generate
/lib/netplan/netplan-dbus
/lib/systemd
/lib/systemd/system
/lib/systemd/system/netplan-wpa@.service
/lib/systemd/system-generators
/usr
/usr/sbin
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/netplan
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/io.netplan.Netplan.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/io.netplan.Netplan.conf
/usr/share/doc
/usr/share/doc/netplan.io
/usr/share/doc/netplan.io/changelog.Debian.gz
/usr/share/doc/netplan.io/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/netplan.5.gz
/usr/share/man/man8
/usr/share/man/man8/netplan-apply.8.gz
/usr/share/man/man8/netplan-generate.8.gz
/usr/share/man/man8/netplan-try.8.gz
/usr/share/netplan
/usr/share/netplan/netplan
/usr/share/netplan/netplan/__init__.py
/usr/share/netplan/netplan/_features.py
/usr/share/netplan/netplan/cli
/usr/share/netplan/netplan/cli/__init__.py
/usr/share/netplan/netplan/cli/commands
/usr/share/netplan/netplan/cli/commands/__init__.py
/usr/share/netplan/netplan/cli/commands/apply.py
/usr/share/netplan/netplan/cli/commands/generate.py
/usr/share/netplan/netplan/cli/commands/info.py
/usr/share/netplan/netplan/cli/commands/ip.py
/usr/share/netplan/netplan/cli/commands/migrate.py
/usr/share/netplan/netplan/cli/commands/try_command.py
/usr/share/netplan/netplan/cli/core.py
/usr/share/netplan/netplan/cli/utils.py
/usr/share/netplan/netplan/configmanager.py
/usr/share/netplan/netplan/terminal.py
/usr/share/netplan/netplan.script
/lib/systemd/system-generators/netplan
/usr/sbin/netplan
