# squashfs-root/var/lib/dpkg/info/network-manager.conffiles
/etc/NetworkManager/NetworkManager.conf
/etc/NetworkManager/conf.d/default-wifi-powersave-on.conf
/etc/NetworkManager/dispatcher.d/01-ifupdown
/etc/dbus-1/system.d/nm-dispatcher.conf
/etc/dbus-1/system.d/org.freedesktop.NetworkManager.conf
/etc/init.d/network-manager
