# squashfs-root/var/lib/dpkg/info/network-manager-config-connectivity-ubuntu.postinst
#!/bin/sh

set -e

if [ "$1" = configure ]; then
    invoke-rc.d network-manager force-reload || true
fi


