# squashfs-root/var/lib/dpkg/info/networkd-dispatcher.list
/.
/etc
/etc/default
/etc/default/networkd-dispatcher
/etc/networkd-dispatcher
/etc/networkd-dispatcher/dormant.d
/etc/networkd-dispatcher/no-carrier.d
/etc/networkd-dispatcher/off.d
/etc/networkd-dispatcher/routable.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/networkd-dispatcher.service
/usr
/usr/bin
/usr/bin/networkd-dispatcher
/usr/lib
/usr/lib/networkd-dispatcher
/usr/lib/networkd-dispatcher/dormant.d
/usr/lib/networkd-dispatcher/no-carrier.d
/usr/lib/networkd-dispatcher/off.d
/usr/lib/networkd-dispatcher/routable.d
/usr/share
/usr/share/doc
/usr/share/doc/networkd-dispatcher
/usr/share/doc/networkd-dispatcher/README.md.gz
/usr/share/doc/networkd-dispatcher/changelog.Debian.gz
/usr/share/doc/networkd-dispatcher/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/networkd-dispatcher.8.gz
