# squashfs-root/var/lib/dpkg/info/packagekit.list
/.
/etc
/etc/PackageKit
/etc/PackageKit/PackageKit.conf
/etc/PackageKit/Vendor.conf
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20packagekit
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.PackageKit.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/packagekit-offline-update.service
/lib/systemd/system/packagekit.service
/lib/systemd/system/system-update.target.wants
/usr
/usr/lib
/usr/lib/packagekit
/usr/lib/packagekit/packagekit-direct
/usr/lib/packagekit/packagekitd
/usr/lib/packagekit/pk-offline-update
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/packagekit-backend
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_aptcc.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_dummy.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_test_fail.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_test_nop.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_test_spawn.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_test_succeed.so
/usr/lib/x86_64-linux-gnu/packagekit-backend/libpk_backend_test_thread.so
/usr/share
/usr/share/PackageKit
/usr/share/PackageKit/helpers
/usr/share/PackageKit/helpers/aptcc
/usr/share/PackageKit/helpers/aptcc/pkconffile
/usr/share/PackageKit/helpers/aptcc/pkconffile.nodiff
/usr/share/PackageKit/helpers/test_spawn
/usr/share/PackageKit/helpers/test_spawn/search-name.sh
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.PackageKit.Transaction.xml
/usr/share/dbus-1/interfaces/org.freedesktop.PackageKit.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.PackageKit.service
/usr/share/doc
/usr/share/doc/packagekit
/usr/share/doc/packagekit/NEWS.gz
/usr/share/doc/packagekit/README.Debian
/usr/share/doc/packagekit/copyright
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.packagekit.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.packagekit.rules
/var
/var/cache
/var/cache/PackageKit
/var/cache/PackageKit/downloads
/var/lib
/var/lib/PackageKit
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/org.freedesktop.packagekit.pkla
/lib/systemd/system/system-update.target.wants/packagekit-offline-update.service
/usr/share/doc/packagekit/changelog.Debian.gz
