# squashfs-root/var/lib/dpkg/info/policykit-1.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.PolicyKit1.conf
/etc/pam.d
/etc/pam.d/polkit-1
/etc/polkit-1
/etc/polkit-1/localauthority
/etc/polkit-1/localauthority/10-vendor.d
/etc/polkit-1/localauthority/20-org.d
/etc/polkit-1/localauthority/30-site.d
/etc/polkit-1/localauthority/50-local.d
/etc/polkit-1/localauthority/90-mandatory.d
/etc/polkit-1/localauthority.conf.d
/etc/polkit-1/localauthority.conf.d/50-localauthority.conf
/etc/polkit-1/localauthority.conf.d/51-ubuntu-admin.conf
/etc/polkit-1/nullbackend.conf.d
/etc/polkit-1/nullbackend.conf.d/50-nullbackend.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/polkit.service
/usr
/usr/bin
/usr/bin/pkaction
/usr/bin/pkcheck
/usr/bin/pkexec
/usr/bin/pkttyagent
/usr/lib
/usr/lib/policykit-1
/usr/lib/policykit-1/polkit-agent-helper-1
/usr/lib/policykit-1/polkitd
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/polkit-1
/usr/lib/x86_64-linux-gnu/polkit-1/extensions
/usr/lib/x86_64-linux-gnu/polkit-1/extensions/libnullbackend.so
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.PolicyKit1.service
/usr/share/doc
/usr/share/doc/policykit-1
/usr/share/doc/policykit-1/NEWS.gz
/usr/share/doc/policykit-1/README
/usr/share/doc/policykit-1/copyright
/usr/share/gettext
/usr/share/gettext/its
/usr/share/gettext/its/polkit.its
/usr/share/gettext/its/polkit.loc
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pkaction.1.gz
/usr/share/man/man1/pkcheck.1.gz
/usr/share/man/man1/pkexec.1.gz
/usr/share/man/man1/pkttyagent.1.gz
/usr/share/man/man8
/usr/share/man/man8/pklocalauthority.8.gz
/usr/share/man/man8/polkit.8.gz
/usr/share/man/man8/polkitd.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.policykit.policy
/var
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/20-org.d
/var/lib/polkit-1/localauthority/30-site.d
/var/lib/polkit-1/localauthority/50-local.d
/var/lib/polkit-1/localauthority/90-mandatory.d
/usr/share/doc/policykit-1/changelog.Debian.gz
