# squashfs-root/var/lib/dpkg/info/popularity-contest.list
/etc
/etc/cron.daily
/etc/cron.daily/popularity-contest
