# squashfs-root/var/lib/dpkg/info/ppp.list
/.
/etc
/etc/chatscripts
/etc/chatscripts/gprs
/etc/chatscripts/pap
/etc/init.d
/etc/init.d/pppd-dns
/etc/logrotate.d
/etc/logrotate.d/ppp
/etc/pam.d
/etc/pam.d/ppp
/etc/ppp
/etc/ppp/ip-down
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/0000usepeerdns
/etc/ppp/ip-up
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/0000usepeerdns
/etc/ppp/ipv6-down
/etc/ppp/ipv6-down.d
/etc/ppp/ipv6-up
/etc/ppp/ipv6-up.d
/etc/ppp/options
/etc/ppp/peers
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/pppd-dns.service
/usr
/usr/bin
/usr/bin/plog
/usr/bin/poff
/usr/bin/pon
/usr/lib
/usr/lib/pppd
/usr/lib/pppd/2.4.7
/usr/lib/pppd/2.4.7/minconn.so
/usr/lib/pppd/2.4.7/openl2tp.so
/usr/lib/pppd/2.4.7/passprompt.so
/usr/lib/pppd/2.4.7/passwordfd.so
/usr/lib/pppd/2.4.7/pppoatm.so
/usr/lib/pppd/2.4.7/pppol2tp.so
/usr/lib/pppd/2.4.7/radattr.so
/usr/lib/pppd/2.4.7/radius.so
/usr/lib/pppd/2.4.7/radrealms.so
/usr/lib/pppd/2.4.7/rp-pppoe.so
/usr/lib/pppd/2.4.7/winbind.so
/usr/sbin
/usr/sbin/chat
/usr/sbin/pppd
/usr/sbin/pppdump
/usr/sbin/pppoe-discovery
/usr/sbin/pppstats
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/poff
/usr/share/bash-completion/completions/pon
/usr/share/doc
/usr/share/doc/ppp
/usr/share/doc/ppp/FAQ.gz
/usr/share/doc/ppp/PLUGINS.gz
/usr/share/doc/ppp/README.Debian.gz
/usr/share/doc/ppp/README.MPPE
/usr/share/doc/ppp/README.MSCHAP80.gz
/usr/share/doc/ppp/README.MSCHAP81
/usr/share/doc/ppp/README.STATIC-IP
/usr/share/doc/ppp/README.cbcp
/usr/share/doc/ppp/README.gz
/usr/share/doc/ppp/README.pppoe
/usr/share/doc/ppp/README.pppol2tp
/usr/share/doc/ppp/README.pwfd
/usr/share/doc/ppp/README.win9x
/usr/share/doc/ppp/SETUP.gz
/usr/share/doc/ppp/changelog.Debian.gz
/usr/share/doc/ppp/copyright
/usr/share/doc/ppp/defaultroute
/usr/share/doc/ppp/examples
/usr/share/doc/ppp/examples/filters
/usr/share/doc/ppp/examples/interfaces
/usr/share/doc/ppp/examples/options.ttyXX
/usr/share/doc/ppp/examples/peers-gprs
/usr/share/doc/ppp/examples/peers-pppoa
/usr/share/doc/ppp/examples/peers-pppoe
/usr/share/doc/ppp/examples/per-linkname
/usr/share/doc/ppp/examples/popp
/usr/share/doc/ppp/examples/scripts
/usr/share/doc/ppp/examples/scripts/README.gz
/usr/share/doc/ppp/examples/scripts/autopppd.gz
/usr/share/doc/ppp/examples/scripts/callback
/usr/share/doc/ppp/examples/scripts/chat-callback
/usr/share/doc/ppp/examples/scripts/chatchat
/usr/share/doc/ppp/examples/scripts/chatchat/README.gz
/usr/share/doc/ppp/examples/scripts/chatchat/chatchat.c.gz
/usr/share/doc/ppp/examples/scripts/ip-down.local.add
/usr/share/doc/ppp/examples/scripts/ip-up.local.add
/usr/share/doc/ppp/examples/scripts/ipv6-down.sample
/usr/share/doc/ppp/examples/scripts/ipv6-up.sample
/usr/share/doc/ppp/examples/scripts/options-rsh-loc
/usr/share/doc/ppp/examples/scripts/options-rsh-rem
/usr/share/doc/ppp/examples/scripts/options-ssh-loc
/usr/share/doc/ppp/examples/scripts/options-ssh-rem
/usr/share/doc/ppp/examples/scripts/plog
/usr/share/doc/ppp/examples/scripts/poff
/usr/share/doc/ppp/examples/scripts/pon
/usr/share/doc/ppp/examples/scripts/pon.1
/usr/share/doc/ppp/examples/scripts/ppp-off
/usr/share/doc/ppp/examples/scripts/ppp-on
/usr/share/doc/ppp/examples/scripts/ppp-on-dialer
/usr/share/doc/ppp/examples/scripts/ppp-on-rsh
/usr/share/doc/ppp/examples/scripts/ppp-on-ssh
/usr/share/doc/ppp/examples/scripts/redialer
/usr/share/doc/ppp/examples/scripts/secure-card
/usr/share/doc/ppp/examples/userscripts-down
/usr/share/doc/ppp/examples/userscripts-up
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/ppp
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pon.1.gz
/usr/share/man/man8
/usr/share/man/man8/chat.8.gz
/usr/share/man/man8/pppd-radattr.8.gz
/usr/share/man/man8/pppd-radius.8.gz
/usr/share/man/man8/pppd.8.gz
/usr/share/man/man8/pppdump.8.gz
/usr/share/man/man8/pppoe-discovery.8.gz
/usr/share/man/man8/pppstats.8.gz
/usr/share/ppp
/usr/share/ppp/chap-secrets
/usr/share/ppp/pap-secrets
/usr/share/ppp/provider.chatscript
/usr/share/ppp/provider.peer
/usr/share/man/man1/plog.1.gz
/usr/share/man/man1/poff.1.gz
