# squashfs-root/var/lib/dpkg/info/rsync.list
/.
/etc
/etc/default
/etc/default/rsync
/etc/init.d
/etc/init.d/rsync
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rsync.service
/usr
/usr/bin
/usr/bin/rsync
/usr/share
/usr/share/doc
/usr/share/doc/rsync
/usr/share/doc/rsync/README.Debian.gz
/usr/share/doc/rsync/README.gz
/usr/share/doc/rsync/TODO.gz
/usr/share/doc/rsync/changelog.Debian.gz
/usr/share/doc/rsync/copyright
/usr/share/doc/rsync/examples
/usr/share/doc/rsync/examples/logrotate.conf.rsync
/usr/share/doc/rsync/examples/rsyncd.conf
/usr/share/doc/rsync/scripts
/usr/share/doc/rsync/scripts/atomic-rsync.gz
/usr/share/doc/rsync/scripts/cull_options.gz
/usr/share/doc/rsync/scripts/cvs2includes.gz
/usr/share/doc/rsync/scripts/file-attr-restore.gz
/usr/share/doc/rsync/scripts/files-to-excludes.gz
/usr/share/doc/rsync/scripts/git-set-file-times.gz
/usr/share/doc/rsync/scripts/logfilter.gz
/usr/share/doc/rsync/scripts/lsh.gz
/usr/share/doc/rsync/scripts/mnt-excl.gz
/usr/share/doc/rsync/scripts/munge-symlinks.gz
/usr/share/doc/rsync/scripts/rrsync.gz
/usr/share/doc/rsync/scripts/rsyncstats.gz
/usr/share/doc/rsync/tech_report.tex.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rsync
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/rsync.1.gz
/usr/share/man/man5
/usr/share/man/man5/rsyncd.conf.5.gz
