# squashfs-root/var/lib/dpkg/info/rsyslog.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/disable
/etc/apparmor.d/force-complain
/etc/apparmor.d/local
/etc/apparmor.d/usr.sbin.rsyslogd
/etc/default
/etc/default/rsyslog
/etc/init.d
/etc/init.d/rsyslog
/etc/logcheck
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/rsyslog
/etc/logrotate.d
/etc/logrotate.d/rsyslog
/etc/rsyslog.conf
/etc/rsyslog.d
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rsyslog.service
/usr
/usr/lib
/usr/lib/rsyslog
/usr/lib/rsyslog/rsyslog-rotate
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/00rsyslog.conf
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/rsyslog
/usr/lib/x86_64-linux-gnu/rsyslog/imfile.so
/usr/lib/x86_64-linux-gnu/rsyslog/imjournal.so
/usr/lib/x86_64-linux-gnu/rsyslog/imklog.so
/usr/lib/x86_64-linux-gnu/rsyslog/imkmsg.so
/usr/lib/x86_64-linux-gnu/rsyslog/immark.so
/usr/lib/x86_64-linux-gnu/rsyslog/impstats.so
/usr/lib/x86_64-linux-gnu/rsyslog/imptcp.so
/usr/lib/x86_64-linux-gnu/rsyslog/imtcp.so
/usr/lib/x86_64-linux-gnu/rsyslog/imudp.so
/usr/lib/x86_64-linux-gnu/rsyslog/imuxsock.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmnet.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmnetstrms.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmnsd_ptcp.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmregexp.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmstrmsrv.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmtcpclt.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmtcpsrv.so
/usr/lib/x86_64-linux-gnu/rsyslog/lmzlibw.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmanon.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmexternal.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmfields.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmjsonparse.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmpstrucdata.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmsequence.so
/usr/lib/x86_64-linux-gnu/rsyslog/mmutf8fix.so
/usr/lib/x86_64-linux-gnu/rsyslog/omjournal.so
/usr/lib/x86_64-linux-gnu/rsyslog/ommail.so
/usr/lib/x86_64-linux-gnu/rsyslog/omprog.so
/usr/lib/x86_64-linux-gnu/rsyslog/omuxsock.so
/usr/lib/x86_64-linux-gnu/rsyslog/pmaixforwardedfrom.so
/usr/lib/x86_64-linux-gnu/rsyslog/pmcisconames.so
/usr/lib/x86_64-linux-gnu/rsyslog/pmlastmsg.so
/usr/lib/x86_64-linux-gnu/rsyslog/pmsnare.so
/usr/sbin
/usr/sbin/rsyslogd
/usr/share
/usr/share/doc
/usr/share/doc/rsyslog
/usr/share/doc/rsyslog/AUTHORS
/usr/share/doc/rsyslog/NEWS.Debian.gz
/usr/share/doc/rsyslog/README.Debian
/usr/share/doc/rsyslog/changelog.Debian.gz
/usr/share/doc/rsyslog/copyright
/usr/share/doc/rsyslog/examples
/usr/share/doc/rsyslog/examples/rsyslog.d
/usr/share/doc/rsyslog/examples/rsyslog.d/console.conf
/usr/share/doc/rsyslog/examples/rsyslog.d/xconsole.conf
/usr/share/doc/rsyslog/examples/tmpfiles.d
/usr/share/doc/rsyslog/examples/tmpfiles.d/xconsole.conf
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/rsyslog
/usr/share/man
/usr/share/man/man1
/usr/share/man/man5
/usr/share/man/man5/rsyslog.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/rsyslogd.8.gz
/usr/share/rsyslog
/usr/share/rsyslog/50-default.conf
/var
/var/spool
/var/spool/rsyslog
