# squashfs-root/var/lib/dpkg/info/rtkit.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.RealtimeKit1.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/rtkit-daemon.service
/usr
/usr/lib
/usr/lib/rtkit
/usr/lib/rtkit/rtkit-daemon
/usr/lib/rtkit/rtkit-test
/usr/sbin
/usr/sbin/rtkitctl
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.RealtimeKit1.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.RealtimeKit1.service
/usr/share/doc
/usr/share/doc/rtkit
/usr/share/doc/rtkit/README.gz
/usr/share/doc/rtkit/changelog.Debian.gz
/usr/share/doc/rtkit/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/rtkitctl.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.RealtimeKit1.policy
