# squashfs-root/var/lib/dpkg/info/sane-utils.conffiles
/etc/default/saned
/etc/init.d/saned
/etc/sane.d/saned.conf
