# squashfs-root/var/lib/dpkg/info/snapd.list
/.
/etc
/etc/apparmor.d
/etc/apparmor.d/usr.lib.snapd.snap-confine.real
/etc/apt
/etc/apt/apt.conf.d
/etc/apt/apt.conf.d/20snapd.conf
/etc/profile.d
/etc/profile.d/apps-bin-path.sh
/etc/xdg
/etc/xdg/autostart
/etc/xdg/autostart/snap-userd-autostart.desktop
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/snapd.autoimport.service
/lib/systemd/system/snapd.core-fixup.service
/lib/systemd/system/snapd.failure.service
/lib/systemd/system/snapd.seeded.service
/lib/systemd/system/snapd.service
/lib/systemd/system/snapd.snap-repair.service
/lib/systemd/system/snapd.snap-repair.timer
/lib/systemd/system/snapd.socket
/lib/systemd/system/snapd.system-shutdown.service
/lib/systemd/system-generators
/lib/systemd/system-generators/snapd-generator
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/66-snapd-autoimport.rules
/snap
/usr
/usr/bin
/usr/bin/snap
/usr/bin/snapfuse
/usr/lib
/usr/lib/environment.d
/usr/lib/environment.d/990-snapd.conf
/usr/lib/snapd
/usr/lib/snapd/complete.sh
/usr/lib/snapd/etelpmoc.sh
/usr/lib/snapd/info
/usr/lib/snapd/snap-confine
/usr/lib/snapd/snap-device-helper
/usr/lib/snapd/snap-discard-ns
/usr/lib/snapd/snap-exec
/usr/lib/snapd/snap-failure
/usr/lib/snapd/snap-gdb-shim
/usr/lib/snapd/snap-mgmt
/usr/lib/snapd/snap-repair
/usr/lib/snapd/snap-seccomp
/usr/lib/snapd/snap-update-ns
/usr/lib/snapd/snapctl
/usr/lib/snapd/snapd
/usr/lib/snapd/snapd.core-fixup.sh
/usr/lib/snapd/snapd.run-from-snap
/usr/lib/snapd/system-shutdown
/usr/lib/systemd
/usr/lib/systemd/system-environment-generators
/usr/lib/systemd/system-environment-generators/snapd-env-generator
/usr/lib/systemd/user
/usr/lib/systemd/user/snapd.session-agent.service
/usr/lib/systemd/user/snapd.session-agent.socket
/usr/lib/systemd/user/sockets.target.wants
/usr/share
/usr/share/applications
/usr/share/applications/snap-handle-link.desktop
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/snap
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/services/io.snapcraft.Launcher.service
/usr/share/dbus-1/services/io.snapcraft.Settings.service
/usr/share/doc
/usr/share/doc/snapd
/usr/share/doc/snapd/changelog.gz
/usr/share/doc/snapd/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/snap-confine.8.gz
/usr/share/man/man8/snap-discard-ns.8.gz
/usr/share/man/man8/snap.8.gz
/usr/share/man/man8/snapd-env-generator.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/io.snapcraft.snapd.policy
/var
/var/cache
/var/cache/snapd
/var/lib
/var/lib/snapd
/var/lib/snapd/apparmor
/var/lib/snapd/apparmor/snap-confine
/var/lib/snapd/auto-import
/var/lib/snapd/desktop
/var/lib/snapd/environment
/var/lib/snapd/firstboot
/var/lib/snapd/lib
/var/lib/snapd/lib/gl
/var/lib/snapd/lib/gl32
/var/lib/snapd/lib/glvnd
/var/lib/snapd/lib/vulkan
/var/lib/snapd/snaps
/var/lib/snapd/snaps/partial
/var/lib/snapd/void
/var/snap
/lib/udev/snappy-app-dev
/usr/bin/snapctl
/usr/bin/ubuntu-core-launcher
/usr/lib/systemd/user/sockets.target.wants/snapd.session-agent.socket
