# squashfs-root/var/lib/dpkg/info/speech-dispatcher.list
/.
/etc
/etc/default
/etc/default/speech-dispatcher
/etc/init.d
/etc/init.d/speech-dispatcher
/etc/logrotate.d
/etc/logrotate.d/speech-dispatcher
/etc/speech-dispatcher
/etc/speech-dispatcher/clients
/etc/speech-dispatcher/clients/emacs.conf
/etc/speech-dispatcher/modules
/etc/speech-dispatcher/modules/cicero.conf
/etc/speech-dispatcher/modules/dtk-generic.conf
/etc/speech-dispatcher/modules/epos-generic.conf
/etc/speech-dispatcher/modules/espeak-generic.conf
/etc/speech-dispatcher/modules/espeak-mbrola-generic.conf
/etc/speech-dispatcher/modules/espeak-ng.conf
/etc/speech-dispatcher/modules/espeak.conf
/etc/speech-dispatcher/modules/festival.conf
/etc/speech-dispatcher/modules/flite.conf
/etc/speech-dispatcher/modules/ibmtts.conf
/etc/speech-dispatcher/modules/ivona.conf
/etc/speech-dispatcher/modules/llia_phon-generic.conf
/etc/speech-dispatcher/modules/pico-generic.conf
/etc/speech-dispatcher/modules/swift-generic.conf
/etc/speech-dispatcher/speechd.conf
/usr
/usr/bin
/usr/bin/spd-say
/usr/bin/speech-dispatcher
/usr/lib
/usr/lib/speech-dispatcher-modules
/usr/lib/speech-dispatcher-modules/sd_dummy
/usr/lib/speech-dispatcher-modules/sd_generic
/usr/share
/usr/share/doc
/usr/share/doc/speech-dispatcher
/usr/share/doc/speech-dispatcher/AUTHORS
/usr/share/doc/speech-dispatcher/NEWS.gz
/usr/share/doc/speech-dispatcher/README
/usr/share/doc/speech-dispatcher/README.Debian
/usr/share/doc/speech-dispatcher/architecture-future.png
/usr/share/doc/speech-dispatcher/architecture.png
/usr/share/doc/speech-dispatcher/copyright
/usr/share/doc/speech-dispatcher/priorities.png
/usr/share/info
/usr/share/info/spd-say.info.gz
/usr/share/info/speech-dispatcher.info.gz
/usr/share/info/ssip.info.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/speech-dispatcher
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/spd-say.1.gz
/usr/share/man/man1/speech-dispatcher.1.gz
/usr/share/sounds
/usr/share/sounds/speech-dispatcher
/usr/share/sounds/speech-dispatcher/dummy-message.wav
/usr/share/speech-dispatcher
/usr/share/speech-dispatcher/conf
/usr/share/speech-dispatcher/conf/clients
/usr/share/speech-dispatcher/conf/clients/emacs.conf
/usr/share/speech-dispatcher/conf/desktop
/usr/share/speech-dispatcher/conf/desktop/speechd.desktop
/usr/share/speech-dispatcher/conf/modules
/usr/share/speech-dispatcher/conf/modules/cicero.conf
/usr/share/speech-dispatcher/conf/modules/espeak-ng.conf
/usr/share/speech-dispatcher/conf/modules/espeak.conf
/usr/share/speech-dispatcher/conf/modules/festival.conf
/usr/share/speech-dispatcher/conf/speechd.conf
/var
/var/log
/var/log/speech-dispatcher
/usr/share/doc/speech-dispatcher/changelog.Debian.gz
