# squashfs-root/var/lib/dpkg/info/spice-vdagent.list
/.
/etc
/etc/init.d
/etc/init.d/spice-vdagent
/etc/xdg
/etc/xdg/autostart
/etc/xdg/autostart/spice-vdagent.desktop
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/spice-vdagentd.service
/lib/systemd/system/spice-vdagentd.target
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/70-spice-vdagentd.rules
/usr
/usr/bin
/usr/bin/spice-vdagent
/usr/lib
/usr/lib/tmpfiles.d
/usr/lib/tmpfiles.d/spice-vdagentd.conf
/usr/sbin
/usr/sbin/spice-vdagentd
/usr/share
/usr/share/doc
/usr/share/doc/spice-vdagent
/usr/share/doc/spice-vdagent/README
/usr/share/doc/spice-vdagent/TODO
/usr/share/doc/spice-vdagent/changelog.Debian.gz
/usr/share/doc/spice-vdagent/copyright
/usr/share/gdm
/usr/share/gdm/autostart
/usr/share/gdm/autostart/LoginWindow
/usr/share/gdm/autostart/LoginWindow/spice-vdagent.desktop
/usr/share/gdm/greeter
/usr/share/gdm/greeter/autostart
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/spice-vdagent.1.gz
/usr/share/man/man1/spice-vdagentd.1.gz
/lib/systemd/system/spice-vdagent.service
/usr/share/gdm/greeter/autostart/spice-vdagent.desktop
