# squashfs-root/var/lib/dpkg/info/spice-vdagent.conffiles
/etc/init.d/spice-vdagent
/etc/xdg/autostart/spice-vdagent.desktop
