# squashfs-root/var/lib/dpkg/info/thermald.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.thermald.conf
/etc/thermald
/etc/thermald/thermal-conf.xml
/etc/thermald/thermal-cpu-cdev-order.xml
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/thermald.service
/usr
/usr/sbin
/usr/sbin/thermald
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.thermald.service
/usr/share/doc
/usr/share/doc/thermald
/usr/share/doc/thermald/changelog.Debian.gz
/usr/share/doc/thermald/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/thermal-conf.xml.5.gz
/usr/share/man/man8
/usr/share/man/man8/thermald.8.gz
