# squashfs-root/var/lib/dpkg/info/ubuntu-advantage-tools.conffiles
/etc/cron.daily/ubuntu-advantage-tools
/etc/update-motd.d/80-esm
/etc/update-motd.d/80-livepatch
