# squashfs-root/var/lib/dpkg/info/ubuntu-drivers-common.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/gpu-manager.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/71-u-d-c-gpu-detection.rules
/sbin
/sbin/u-d-c-print-pci-ids
/usr
/usr/bin
/usr/bin/gpu-manager
/usr/bin/nvidia-detector
/usr/bin/quirks-handler
/usr/bin/ubuntu-drivers
/usr/lib
/usr/lib/nvidia
/usr/lib/nvidia/pre-install
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/NvidiaDetector
/usr/lib/python3/dist-packages/NvidiaDetector/__init__.py
/usr/lib/python3/dist-packages/NvidiaDetector/alternatives.py
/usr/lib/python3/dist-packages/NvidiaDetector/nvidiadetector.py
/usr/lib/python3/dist-packages/Quirks
/usr/lib/python3/dist-packages/Quirks/__init__.py
/usr/lib/python3/dist-packages/Quirks/quirkapplier.py
/usr/lib/python3/dist-packages/Quirks/quirkinfo.py
/usr/lib/python3/dist-packages/Quirks/quirkreader.py
/usr/lib/python3/dist-packages/UbuntuDrivers
/usr/lib/python3/dist-packages/UbuntuDrivers/__init__.py
/usr/lib/python3/dist-packages/UbuntuDrivers/detect.py
/usr/lib/python3/dist-packages/UbuntuDrivers/kerneldetection.py
/usr/lib/python3/dist-packages/ubuntu_drivers_common-0.0.0.egg-info
/usr/lib/python3/dist-packages/ubuntu_drivers_common-0.0.0.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/ubuntu_drivers_common-0.0.0.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/ubuntu_drivers_common-0.0.0.egg-info/top_level.txt
/usr/lib/ubiquity
/usr/lib/ubiquity/target-config
/usr/lib/ubiquity/target-config/31ubuntu_driver_packages
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/ubuntu-drivers-common.py
/usr/share/doc
/usr/share/doc/ubuntu-drivers-common
/usr/share/doc/ubuntu-drivers-common/README
/usr/share/doc/ubuntu-drivers-common/changelog.gz
/usr/share/doc/ubuntu-drivers-common/copyright
/usr/share/python3
/usr/share/python3/runtime.d
/usr/share/python3/runtime.d/ubuntu-drivers-common.rtupdate
/usr/share/ubuntu-drivers-common
/usr/share/ubuntu-drivers-common/detect
/usr/share/ubuntu-drivers-common/detect/arm-gles.py
/usr/share/ubuntu-drivers-common/detect/open-vm-dkms.py
/usr/share/ubuntu-drivers-common/detect/sl-modem.py
/usr/share/ubuntu-drivers-common/fake-devices-wrapper
/usr/share/ubuntu-drivers-common/obsolete
/usr/share/ubuntu-drivers-common/quirks
/usr/share/ubuntu-drivers-common/quirks/dell_latitude
/usr/share/ubuntu-drivers-common/quirks/lenovo_thinkpad
/usr/share/ubuntu-drivers-common/quirks/put_your_quirks_here
/var
/var/lib
/var/lib/ubuntu-drivers-common
