# squashfs-root/var/lib/dpkg/info/ubuntu-advantage-tools.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/ubuntu-advantage-tools
/etc/update-motd.d
/etc/update-motd.d/80-esm
/etc/update-motd.d/80-livepatch
/usr
/usr/bin
/usr/bin/ubuntu-advantage
/usr/share
/usr/share/doc
/usr/share/doc/ubuntu-advantage-tools
/usr/share/doc/ubuntu-advantage-tools/changelog.gz
/usr/share/doc/ubuntu-advantage-tools/copyright
/usr/share/keyrings
/usr/share/keyrings/ubuntu-esm-keyring.gpg
/usr/share/keyrings/ubuntu-fips-keyring.gpg
/usr/share/keyrings/ubuntu-fips-updates-keyring.gpg
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/ubuntu-advantage.1.gz
/usr/share/ubuntu-advantage-tools
/usr/share/ubuntu-advantage-tools/modules
/usr/share/ubuntu-advantage-tools/modules/apt.sh
/usr/share/ubuntu-advantage-tools/modules/credentials.sh
/usr/share/ubuntu-advantage-tools/modules/hardware.sh
/usr/share/ubuntu-advantage-tools/modules/service-esm.sh
/usr/share/ubuntu-advantage-tools/modules/service-fips.sh
/usr/share/ubuntu-advantage-tools/modules/service-livepatch.sh
/usr/share/ubuntu-advantage-tools/modules/service.sh
/usr/share/ubuntu-advantage-tools/modules/snap.sh
/usr/share/ubuntu-advantage-tools/modules/support.sh
/usr/share/ubuntu-advantage-tools/modules/utils.sh
/var
/var/cache
/var/cache/ubuntu-advantage-tools
/usr/bin/ua
/usr/share/man/man1/ua.1.gz
