# squashfs-root/var/lib/dpkg/info/udisks2.list
/.
/etc
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.UDisks2.conf
/etc/udisks2
/etc/udisks2/udisks2.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/clean-mount-point@.service
/lib/systemd/system/udisks2.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/80-udisks2.rules
/sbin
/sbin/umount.udisks2
/usr
/usr/bin
/usr/bin/udisksctl
/usr/lib
/usr/lib/udisks2
/usr/lib/udisks2/udisks2-inhibit
/usr/lib/udisks2/udisksd
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/udisks2.py
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/udisksctl
/usr/share/dbus-1
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.UDisks2.service
/usr/share/doc
/usr/share/doc/udisks2
/usr/share/doc/udisks2/AUTHORS
/usr/share/doc/udisks2/NEWS.gz
/usr/share/doc/udisks2/README.md
/usr/share/doc/udisks2/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/udisks2
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/udisksctl.1.gz
/usr/share/man/man5
/usr/share/man/man5/udisks2.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/udisks.8.gz
/usr/share/man/man8/udisksd.8.gz
/usr/share/man/man8/umount.udisks2.8.gz
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.UDisks2.policy
/usr/share/doc/udisks2/changelog.Debian.gz
