# squashfs-root/var/lib/dpkg/info/ufw.conffiles
/etc/default/ufw
/etc/init.d/ufw
/etc/logrotate.d/ufw
/etc/rsyslog.d/20-ufw.conf
/etc/ufw/sysctl.conf
