# squashfs-root/var/lib/dpkg/info/unattended-upgrades.list
/.
/etc
/etc/apt
/etc/apt/apt.conf.d
/etc/init.d
/etc/init.d/unattended-upgrades
/etc/kernel
/etc/kernel/postinst.d
/etc/kernel/postinst.d/unattended-upgrades
/etc/logrotate.d
/etc/logrotate.d/unattended-upgrades
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_unattended-upgrades-hibernate
/etc/update-motd.d
/etc/update-motd.d/92-unattended-upgrades
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/unattended-upgrades.service
/lib/systemd/system-sleep
/lib/systemd/system-sleep/unattended-upgrades
/usr
/usr/bin
/usr/bin/unattended-upgrade
/usr/lib
/usr/lib/python3
/usr/lib/python3/dist-packages
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/PKG-INFO
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/dependency_links.txt
/usr/lib/python3/dist-packages/unattended_upgrades-0.1.egg-info/top_level.txt
/usr/lib/systemd
/usr/lib/systemd/logind.conf.d
/usr/lib/systemd/logind.conf.d/unattended-upgrades-logind-maxdelay.conf
/usr/share
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_unattended-upgrades.py
/usr/share/doc
/usr/share/doc/unattended-upgrades
/usr/share/doc/unattended-upgrades/NEWS.Debian.gz
/usr/share/doc/unattended-upgrades/README.md.gz
/usr/share/doc/unattended-upgrades/changelog.gz
/usr/share/doc/unattended-upgrades/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/unattended-upgrades
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/unattended-upgrade.8.gz
/usr/share/unattended-upgrades
/usr/share/unattended-upgrades/20auto-upgrades
/usr/share/unattended-upgrades/20auto-upgrades-disabled
/usr/share/unattended-upgrades/50unattended-upgrades
/usr/share/unattended-upgrades/unattended-upgrade-shutdown
/usr/share/unattended-upgrades/update-motd-unattended-upgrades
/var
/var/lib
/var/lib/unattended-upgrades
/var/log
/var/log/unattended-upgrades
/usr/bin/unattended-upgrades
/usr/share/man/man8/unattended-upgrades.8.gz
