# squashfs-root/var/lib/dpkg/info/unattended-upgrades.conffiles
/etc/init.d/unattended-upgrades
/etc/kernel/postinst.d/unattended-upgrades
/etc/logrotate.d/unattended-upgrades
/etc/pm/sleep.d/10_unattended-upgrades-hibernate
/etc/update-motd.d/92-unattended-upgrades
