# squashfs-root/var/lib/dpkg/info/update-notifier-common.conffiles
/etc/apt/apt.conf.d/10periodic
/etc/apt/apt.conf.d/15update-stamp
/etc/apt/apt.conf.d/20archive
/etc/apt/apt.conf.d/99update-notifier
/etc/cron.daily/update-notifier-common
/etc/cron.weekly/update-notifier-common
/etc/update-motd.d/90-updates-available
/etc/update-motd.d/95-hwe-eol
/etc/update-motd.d/98-fsck-at-reboot
/etc/update-motd.d/98-reboot-required
