# squashfs-root/var/lib/dpkg/info/upower.list
/.
/etc
/etc/UPower
/etc/UPower/UPower.conf
/etc/dbus-1
/etc/dbus-1/system.d
/etc/dbus-1/system.d/org.freedesktop.UPower.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/upower.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/95-upower-csr.rules
/lib/udev/rules.d/95-upower-hid.rules
/lib/udev/rules.d/95-upower-wup.rules
/usr
/usr/bin
/usr/bin/upower
/usr/lib
/usr/lib/upower
/usr/lib/upower/upowerd
/usr/share
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.UPower.Device.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UPower.KbdBacklight.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UPower.Wakeups.xml
/usr/share/dbus-1/interfaces/org.freedesktop.UPower.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.UPower.service
/usr/share/doc
/usr/share/doc/upower
/usr/share/doc/upower/AUTHORS
/usr/share/doc/upower/NEWS.gz
/usr/share/doc/upower/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/upower
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/upower.1.gz
/usr/share/man/man7
/usr/share/man/man7/UPower.7.gz
/usr/share/man/man8
/usr/share/man/man8/upowerd.8.gz
/var
/var/lib
/var/lib/upower
/usr/share/doc/upower/changelog.Debian.gz
