# squashfs-root/var/lib/dpkg/info/usb-modeswitch.list
/.
/etc
/etc/usb_modeswitch.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/usb_modeswitch@.service
/lib/udev
/lib/udev/usb_modeswitch
/usr
/usr/sbin
/usr/sbin/usb_modeswitch
/usr/sbin/usb_modeswitch_dispatcher
/usr/share
/usr/share/doc
/usr/share/doc/usb-modeswitch
/usr/share/doc/usb-modeswitch/NEWS.Debian.gz
/usr/share/doc/usb-modeswitch/README.Debian
/usr/share/doc/usb-modeswitch/README.gz
/usr/share/doc/usb-modeswitch/changelog.Debian.gz
/usr/share/doc/usb-modeswitch/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/usb_modeswitch.1.gz
/usr/share/man/man1/usb_modeswitch_dispatcher.1.gz
/var
/var/lib
/var/lib/usb_modeswitch
