# squashfs-root/var/lib/dpkg/info/usbmuxd.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/usbmuxd.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/39-usbmuxd.rules
/usr
/usr/sbin
/usr/sbin/usbmuxd
/usr/share
/usr/share/doc
/usr/share/doc/usbmuxd
/usr/share/doc/usbmuxd/AUTHORS
/usr/share/doc/usbmuxd/README
/usr/share/doc/usbmuxd/changelog.Debian.gz
/usr/share/doc/usbmuxd/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/usbmuxd.1.gz
