# squashfs-root/var/lib/dpkg/info/whoopsie.conffiles
/etc/init.d/whoopsie
/etc/init/whoopsie.conf
