# squashfs-root/var/lib/dpkg/info/xserver-xorg-input-wacom-hwe-18.04.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/wacom-inputattach@.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/69-wacom.rules
/usr
/usr/bin
/usr/bin/isdv4-serial-debugger
/usr/bin/isdv4-serial-inputattach
/usr/bin/xsetwacom
/usr/include
/usr/include/xorg
/usr/include/xorg/Xwacom.h
/usr/include/xorg/isdv4.h
/usr/include/xorg/wacom-properties.h
/usr/include/xorg/wacom-util.h
/usr/lib
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/pkgconfig
/usr/lib/x86_64-linux-gnu/pkgconfig/xorg-wacom.pc
/usr/lib/xorg
/usr/lib/xorg/modules
/usr/lib/xorg/modules/input
/usr/lib/xorg/modules/input/wacom_drv.so
/usr/share
/usr/share/X11
/usr/share/X11/xorg.conf.d
/usr/share/X11/xorg.conf.d/70-wacom.conf
/usr/share/doc
/usr/share/doc/xserver-xorg-input-wacom-hwe-18.04
/usr/share/doc/xserver-xorg-input-wacom-hwe-18.04/changelog.Debian.gz
/usr/share/doc/xserver-xorg-input-wacom-hwe-18.04/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/xsetwacom.1.gz
/usr/share/man/man4
/usr/share/man/man4/wacom.4.gz
