# squashfs-root/snap/core18/1066/lib/systemd/system/user@.service.d/timeout.conf
# Avoid long hangs during shutdown if user services fail/hang due to X.org
# going away too early
[Service]
TimeoutStopSec=5
