# squashfs-root/snap/core/7270/lib/systemd/system/wpa_supplicant.service.d/snap.conf
[Unit]
ConditionPathExists=!/etc/systemd/system/snap.wpa-supplicant.wpa.service
